# roles/gitea_runner/templates/gitea_runner.service.j2
# jinja2: lstrip_blocks: "True"
# gitea_runner.service #######################################################################
[Unit]
Description=gitea_runner
After=network-online.target
Wants=network-online.target
# gitea_runner connects to the rootless podman.socket inside the gitea_runner
# user manager; wait for that user manager to be up before starting.
After=user@{{ gitea_runner_user.uid }}.service
Requires=user@{{ gitea_runner_user.uid }}.service
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
{% endif %}
{% if zfs_has_scratch_mount %}
After=zfs_mount_mnt_scratch.service
Requires=zfs_mount_mnt_scratch.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Restart=always
TimeoutStartSec=120
# gitea-runner holds each in-flight job's container until the workflow
# exits; a tight TimeoutStopSec SIGKILLs mid-job on restart, leaves a
# "runner died" entry on the gitea side, and bricks any queue entry
# waiting on that job. 120s gives capacity=4 jobs a reasonable drain
# window; matches the rationale for github_runner.service.j2's 60s.
TimeoutStopSec=120
RestartSec=5
SyslogIdentifier=%N

Type=simple
User=gitea_runner
Group=gitea_runner
WorkingDirectory=/mnt/services/gitea_runner/
# Wait for the gitea_runner user manager to bring podman.socket up.
# Bounded: 60 × 1s, well inside TimeoutStartSec=120. If the socket
# never appears, the journal records the explicit cause rather than a
# generic start-timeout.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -S /run/user/{{ gitea_runner_user.uid }}/podman/podman.sock ] && exit 0; sleep 1; done; echo "podman.sock not present after 60s" >&2; exit 1'
# Upstream renamed the binary act_runner -> gitea-runner in v1.0.0.
ExecStart=/usr/local/bin/gitea-runner daemon --config /mnt/services/gitea_runner/config.yml
# No ExecReload: gitea-runner v1.0.3's signal handler treats SIGHUP as
# a termination signal, so `systemctl reload` would kill the daemon
# instead of re-reading config. Config changes propagate through a
# full restart, fired by the systemd_unit service helper's restart:
# trigger chain in tasks/main.yml.

[Install]
WantedBy=default.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
