# roles/gitea_runner/templates/gitea_runner.service.j2
#jinja2: lstrip_blocks: "True"
# gitea_runner.service #######################################################################
[Unit]
Description=gitea_runner
After=network-online.target
Wants=network-online.target
# gitea_runner connects to the rootless podman.socket inside the gitea_runner
# user manager; wait for that user manager to be up before starting.
After=user@{{ gitea_runner_user.uid }}.service
Requires=user@{{ gitea_runner_user.uid }}.service
# Three start attempts spaced RestartSec=30s apart = ~90s of wall-clock
# before the burst trips. Tighter spacing (the old RestartSec=5s) tripped
# in ~15s and then blocked auto-restart for the remaining ~585s of the
# interval -- a transient gitea unreachability at host boot would dead
# the daemon for ~10min before the operator could intervene.
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
{% endif %}
{% if zfs_has_scratch_mount %}
After=zfs_mount_mnt_scratch.service
Requires=zfs_mount_mnt_scratch.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Restart=always
# 60s = 1s daemon settle + the ExecStartPre's max 60s wait-loop. Previous
# 120s was copy-paste padding; if the user manager + podman.socket haven't
# come up in 60s the operator wants the failure surfaced, not silently
# waited for.
TimeoutStartSec=60
# gitea-runner holds each in-flight job's container until the workflow
# exits; a tight TimeoutStopSec SIGKILLs mid-job on restart, leaves a
# "runner died" entry on the gitea side, and bricks any queue entry
# waiting on that job. 120s gives capacity=4 jobs a reasonable drain
# window. Composes with runner.shutdown_timeout=110s in config.yml.j2
# (10s headroom for daemon cleanup before SIGKILL).
TimeoutStopSec=120
RestartSec=30
SyslogIdentifier=%N

# The daemon runs natively (not containerised like roles/github_runner)
# because gitea-runner is a static Go binary that only orchestrates job
# containers via the rootless podman socket -- it doesn't execute any
# user code in-process. Wrapping it in a container would add ~25 lines
# of unit config and nested-DooD complexity for no isolation gain that
# User=gitea_runner doesn't already provide.
Type=simple
User=gitea_runner
Group=gitea_runner
WorkingDirectory=/mnt/services/gitea_runner/
# user@<uid>.service uses Type=notify and only reports active after the
# user manager has reached its default target, so in theory podman.socket
# (pulled in via sockets.target.wants/) is already listening by then. In
# practice -- especially during host boot when ordering is loose -- the
# socket inode can exist for a few hundred ms before podman is accepting
# connections, and `[ -S ]` flips true the moment the inode appears.
# This loop is defensive against that boot-time race; cost when not
# racing is one is-it-there check. Bounded: 60 x 1s, inside
# TimeoutStartSec=60. Failure prints an explicit cause instead of a
# generic systemd start-timeout. (require_docker=true in config.yml.j2
# is the daemon-side counterpart that catches socket-exists-but-not-
# serving once the daemon actually tries to dial it.)
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -S /run/user/{{ gitea_runner_user.uid }}/podman/podman.sock ] && exit 0; sleep 1; done; echo "podman.sock not present after 60s" >&2; exit 1'
# Upstream renamed the binary act_runner -> gitea-runner in v1.0.0.
ExecStart=/usr/local/bin/gitea-runner daemon --config /mnt/services/gitea_runner/config.yml
# No ExecReload: gitea-runner v1.0.3's signal handler treats SIGHUP as
# a termination signal, so `systemctl reload` would kill the daemon
# instead of re-reading config. Config changes propagate through a
# full restart, fired by the systemd_unit service helper's restart:
# trigger chain in tasks/main.yml.

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
