# roles/powertop/files/powertop.service
[Unit]
Description=PowerTOP auto tune

[Service]
Type=oneshot
Environment="TERM=dumb"
RemainAfterExit=true
ExecStart=/usr/sbin/powertop --auto-tune
SyslogIdentifier=%N

[Install]
WantedBy=multi-user.target
