# roles/pihole/templates/pihole.service.j2
#jinja2: lstrip_blocks: True
# pihole.service #######################################################################
[Unit]
Description=pihole
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{# pi-hole's upstream is dnscrypt-proxy on the podman bridge gateway:5053  #}
{# (FTLCONF_dns_upstreams below). Without this ordering the two units stop #}
{# in parallel at shutdown; dnscrypt-proxy exits in <1s while pi-hole has  #}
{# a 120s container stop-timeout, so pi-hole outlives its upstream and any #}
{# cache miss from a LAN client (host services, peers using the keepalived #}
{# VIP) hangs until pi-hole's drain finishes. Order on the .socket so the  #}
{# daemon stays lazily socket-activated; Wants= (not Requires=) so a       #}
{# broken dnscrypt doesn't block pi-hole from serving cached answers.      #}
After=dnscrypt-proxy.socket
Wants=dnscrypt-proxy.socket
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=120
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name pihole \
    --init \
    --log-driver journald \
    --log-opt tag="pihole" \
    --volume /mnt/services/pihole/etc:/etc/pihole \
    --volume /mnt/services/sqlite/pihole:/var/lib/pihole \
    --volume /var/log/pihole:/var/log/pihole \
    --env TZ={{ system_timezone }} \
    --env FTLCONF_dns_listeningMode=SINGLE \
    --env FTLCONF_dns_reply_host_IPv4={{ ansible_facts.default_ipv4.address }} \
    --env FTLCONF_dns_revServers="true,{{ site_subnets.home }},{{ site_subnets.home | ansible.utils.ipmath(1) }}" \
    --env FTLCONF_dns_upstreams="{{ podman_default_network | ansible.utils.ipmath(1) }}#{{ service_ports.dnscrypt_proxy }}" \
    {# Validate DNSSEC at FTL. dnscrypt-proxy already authenticates the   #}
    {# transport to the upstream; this closes the last hop (upstream->FTL #}
    {# ->clients) for signed zones. Local/split-horizon names and the     #}
    {# conditional-forward (revServers) zone are unsigned and unaffected. #}
    --env FTLCONF_dns_dnssec=true \
    {# Split-horizon LAN shortcuts injected straight into FTL's dnsmasq    #}
    {# config rather than via a mounted /etc/dnsmasq.d file. address=/n/ip #}
    {# matches the name + every subdomain; auth.fahm.fr lives at the apex  #}
    {# (not under .lab.) so it needs its own entry -- public Cloudflare    #}
    {# CNAMEs it to lab.fahm.fr, and LAN clients shortcut to lab's         #}
    {# internal IP to avoid a WAN / wireguard hairpin.                     #}
    --env FTLCONF_misc_dnsmasq_lines="address=/pug.fahm.fr/{{ external_ips.pug }};address=/lab.fahm.fr/{{ external_ips.lab }};address=/box.fahm.fr/{{ external_ips.box }};address=/auth.fahm.fr/{{ external_ips.lab }}" \
    {# Admin UI sits behind authelia forward-auth on nginx (only ingress, #}
    {# since the container binds 127.0.0.1:8943) -- pihole's own password #}
    {# is redundant. `webserver.api.password=""` is the supported "no     #}
    {# auth" toggle (a pseudo-item that clears `pwhash` on disk); setting #}
    {# `_pwhash` directly to "" doesn't take effect.                      #}
    --env FTLCONF_webserver_api_password="" \
    --env FTLCONF_database_maxDBdays=31 \
    --env FTLCONF_files_database=/var/lib/pihole/pihole-FTL.db \
    --env FTLCONF_files_gravity=/var/lib/pihole/gravity.db \
    {# pihole-FTL defaults misc.nice=-10 (renice itself) and ntp.sync.active=true #}
    {# (sync system time from upstream). Neither succeeds without CAP_SYS_NICE / #}
    {# CAP_SYS_TIME, which the container doesn't carry -- each failure warns on  #}
    {# every restart. We don't need either: the host's chrony already handles    #}
    {# time, and pihole runs fine at nice=0.                                     #}
    --env FTLCONF_misc_nice=0 \
    --env FTLCONF_ntp_sync_active=false \
    --env PIHOLE_UID={{ pihole_user.uid }} \
    --env PIHOLE_GID={{ pihole_user.group }} \
    --dns=127.0.0.1 \
    --dns={{ bootstrap_dns_resolver }} \
{# Static /etc/hosts overrides inside the container, sharing the same #}
{# `etc_hosts_entries` map roles/hostname renders into /etc/hosts on the #}
{# test box. Bypasses pihole's own split-horizon wildcard #}
{# (`address=/lab.fahm.fr/...`) which would resolve nexus.lab.fahm.fr to #}
{# the fake test external_ips.lab. Production leaves the map unset → #}
{# loop emits no flags. #}
{% for host, ip in (etc_hosts_entries | default({})).items() %}
    --add-host={{ host }}:{{ ip }} \
{% endfor %}
{# Bind only on the host's LAN IP and (in prod) the keepalived VIP -- #}
{# not on 0.0.0.0 -- so the podman bridge gateway addresses stay free #}
{# for aardvark-dns to listen on. keepalived sets ip_nonlocal_bind=1, #}
{# so the VIP bind succeeds even on the standby host. The qemu test #}
{# fixture has no keepalived, so the VIP isn't on any iface and the #}
{# bind would fail with EADDRNOTAVAIL -- skip it under qemu_test. #}
    --publish {{ ansible_facts.default_ipv4.address }}:53:53/udp \
    --publish {{ ansible_facts.default_ipv4.address }}:53:53/tcp \
{% if not qemu_test and pihole_keepalived_ip | ansible.utils.ipaddr %}
    --publish {{ pihole_keepalived_ip }}:53:53/udp \
    --publish {{ pihole_keepalived_ip }}:53:53/tcp \
{% endif %}
    --publish 127.0.0.1:{{ service_ports.pihole }}:80/tcp \
{# Probe via glibc's nss-dns (libresolv) rather than dig: BIND 9.18+ tools #}
{# use libuv's io_uring backend, which podman's default seccomp profile    #}
{# logs (SCMP_ACT_LOG, syscall=425). Every probe emits 2 audit lines       #}
{# (~6k/day on the journal). getent uses old-school send/recv and is       #}
{# silent. pi.hole isn't in /etc/hosts so nss falls through to             #}
{# nameserver 127.0.0.1 — same query target as before. See                 #}
{# containers/podman#18904 for the upstream seccomp story.                 #}
    --health-cmd         '["getent","hosts","pi.hole"]' \
    --health-startup-cmd '["getent","hosts","pi.hole"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 120 \
    docker.io/pihole/pihole:2025.10.0
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
