# roles/pihole/templates/pihole.service.j2
# jinja2: lstrip_blocks: "True"
# pihole.service #######################################################################
[Unit]
Description=pihole
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{# pi-hole's upstream is dnscrypt-proxy on the podman bridge gateway:5053  #}
{# (FTLCONF_dns_upstreams below). Without this ordering the two units stop #}
{# in parallel at shutdown; dnscrypt-proxy exits in <1s while pi-hole has  #}
{# a 120s container stop-timeout, so pi-hole outlives its upstream and any #}
{# cache miss from a LAN client (host services, peers using the keepalived #}
{# VIP) hangs until pi-hole's drain finishes. Order on the .socket so the  #}
{# daemon stays lazily socket-activated; Wants= (not Requires=) so a       #}
{# broken dnscrypt doesn't block pi-hole from serving cached answers.      #}
After=dnscrypt-proxy.socket
Wants=dnscrypt-proxy.socket
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=120
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name pihole \
    --log-driver journald \
    --log-opt tag="pihole" \
    --volume /mnt/services/pihole/etc:/etc/pihole \
    --volume /mnt/services/pihole/dnsmasq:/etc/dnsmasq.d \
    --volume /mnt/services/sqlite/pihole:/var/lib/pihole \
    --volume /var/log/pihole:/var/log/pihole \
    --env TZ={{ system_timezone }} \
    --env FTLCONF_dns_listeningMode=SINGLE \
    --env FTLCONF_dns_reply_host_IPv4={{ ansible_facts.default_ipv4.address }} \
    --env FTLCONF_dns_revServers="true,{{ home_subnet }},{{ home_subnet | ansible.utils.ipmath(1) }}" \
    --env FTLCONF_dns_upstreams="{{ podman_default_network | ansible.utils.ipmath(1) }}#{{ service_ports.dnscrypt_proxy }}" \
    --env FTLCONF_misc_etc_dnsmasq_d=true \
    --env FTLCONF_webserver_api_pwhash="{{ pihole_password | hash('sha256') | hash('sha256') }}" \
    --env FTLCONF_database_maxDBdays=31 \
    --env FTLCONF_files_database=/var/lib/pihole/pihole-FTL.db \
    --env FTLCONF_files_gravity=/var/lib/pihole/gravity.db \
    --env PIHOLE_UID={{ pihole_user.uid }} \
    --env PIHOLE_GID={{ pihole_user.group }} \
    --dns=127.0.0.1 \
    --dns=1.1.1.1 \
{# Static /etc/hosts overrides inside the container, sharing the same #}
{# `extra_hosts` map roles/_test/tasks/mirrors.yml uses on the test box. #}
{# Bypasses pihole's own split-horizon wildcard (`address=/lab.fahm.fr/...`) #}
{# which would resolve nexus.lab.fahm.fr to the fake test external_ips.lab. #}
{# Production leaves the map unset → loop emits no flags. #}
{% for host, ip in (extra_hosts | default({})).items() %}
    --add-host={{ host }}:{{ ip }} \
{% endfor %}
{# Bind only on the host's LAN IP and (in prod) the keepalived VIP -- #}
{# not on 0.0.0.0 -- so the podman bridge gateway addresses stay free #}
{# for aardvark-dns to listen on. keepalived sets ip_nonlocal_bind=1, #}
{# so the VIP bind succeeds even on the standby host. The test #}
{# group_vars leave pihole_keepalived_ip as the literal "changeme" #}
{# placeholder; skip the VIP publish there to keep the test passing. #}
    --publish {{ ansible_facts.default_ipv4.address }}:53:53/udp \
    --publish {{ ansible_facts.default_ipv4.address }}:53:53/tcp \
{% if pihole_keepalived_ip | ansible.utils.ipaddr %}
    --publish {{ pihole_keepalived_ip }}:53:53/udp \
    --publish {{ pihole_keepalived_ip }}:53:53/tcp \
{% endif %}
    --publish 127.0.0.1:8943:80/tcp \
{# Probe via glibc's nss-dns (libresolv) rather than dig: BIND 9.18+ tools #}
{# use libuv's io_uring backend, which podman's default seccomp profile    #}
{# logs (SCMP_ACT_LOG, syscall=425). Every probe emits 2 audit lines       #}
{# (~6k/day on the journal). getent uses old-school send/recv and is       #}
{# silent. pi.hole isn't in /etc/hosts so nss falls through to             #}
{# nameserver 127.0.0.1 — same query target as before. See                 #}
{# containers/podman#18904 for the upstream seccomp story.                 #}
    --health-cmd         '["getent","hosts","pi.hole"]' \
    --health-startup-cmd '["getent","hosts","pi.hole"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 120 \
    docker.io/pihole/pihole:2025.10.0
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=default.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
