# roles/seerr/templates/seerr.service.j2
#jinja2: lstrip_blocks: True
# seerr.service ###########################################################################
[Unit]
Description=seerr
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
# Start deadline for --sdnotify=healthy. seerr runs DB migrations on first boot,
# which on a cold start (a fresh EC2 test cell, or a prod host's first boot after
# an image pull) can run past the 120s this used to allow. Generous; Restart=always
# still recovers a genuinely wedged start.
TimeoutStartSec=300
# Exceed the container's --stop-timeout (120s) so a slow graceful shutdown --
# during which seerr writes settings.json back (the slurp-merge idempotence
# below depends on that writeback) -- finishes before systemd's stop deadline
# SIGKILLs it and skips the writeback.
TimeoutStopSec=180
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
# `USER node:node` is baked into the v3 image, so the entrypoint runs as
# in-container uid/gid 1000 from the start (no privilege drop). The base
# 0:0:65536 range maps the userns the usual way (container root -> host root,
# never exercised here since nothing runs as uid 0); the +1000 override then
# routes the baked-in node user to host `seerr` so the bind-mounted config is
# owned by a real host user. There is deliberately no +0 override: the image
# never runs as in-container root, so there is no root step to remap.
# --init runs a PID-1 reaper: seerr (Node) spawns child processes and a bare
# node as PID 1 would leak their zombies.
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --init \
    --detach \
    --replace \
    --rm \
    --name seerr \
    --security-opt=no-new-privileges \
    --log-driver journald \
    --log-opt tag="seerr" \
    --user 1000:1000 \
    {{ seerr_user | podman_idmap_args(container_uid=1000) | join(" ") }} \
    --env TZ={{ system_timezone }} \
    --volume /mnt/services/seerr:/app/config \
    --volume /mnt/services/sqlite/seerr:/app/config/db \
    --network media \
    --publish 127.0.0.1:{{ service_ports.seerr }}:5055/tcp \
    --health-cmd         {{ "http://localhost:5055/api/v1/status" | podman_health_wget }} \
    --health-startup-cmd {{ "http://localhost:5055/api/v1/status" | podman_health_wget }} \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 120 \
    {{ seerr_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
