# roles/homepage/templates/homepage.service.j2
#jinja2: lstrip_blocks: True
# homepage.service ###########################################################
[Unit]
Description=homepage
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
{% endif %}
{% endif %}
# Deliberately NOT PartOf=zfs_autosnapshot.target. The nightly snapshot
# quiesces stateful services, but homepage has no writable state -- its only
# mount is /mnt/services/homepage:/app/config:ro. Quiescing buys nothing.

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=60
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

# HOMEPAGE_ALLOWED_HOSTS — gethomepage's Host-header allowlist. Three entries:
#   - homepage.<host>.<domain>: nginx vhost (browser path)
#   - localhost:<service_ports.homepage>: host-side _verify probe (ansible uri
#     sends Host: localhost:<port> from the URL netloc; podman forwards TCP
#     transparently so the container sees that literal Host)
#   - localhost:3000: in-container --health-cmd wget probe
# Direct-IP / raw inventory_hostname Host headers are intentionally rejected:
# no legitimate path uses them, and a narrow allowlist hardens against
# Host-header smuggling if a future nginx misconfig forwards a client Host:.
ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name homepage \
    --init \
    --log-driver journald \
    --log-opt tag="homepage" \
    --security-opt=no-new-privileges \
    --network homepage \
    --init \
    --user 0:0 \
    {{ homepage_user | podman_idmap_args | join(" ") }} \
    --volume /mnt/services/homepage:/app/config:ro \
    --tz={{ system_timezone }} \
    --env HOMEPAGE_ALLOWED_HOSTS=homepage.{{ inventory_hostname }}.{{ domain }},localhost:{{ service_ports.homepage }},localhost:3000 \
    --publish 127.0.0.1:{{ service_ports.homepage }}:3000/tcp \
    --health-cmd {{ "http://localhost:3000/api/healthcheck" | podman_health_wget }} \
    --health-startup-cmd {{ "http://localhost:3000/api/healthcheck" | podman_health_wget }} \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "30s" \
    --stop-timeout 30 \
    {{ homepage_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
