# roles/homepage/templates/homepage_alerts.service.j2
# {{ ansible_managed }}
[Unit]
Description=homepage_alerts
Wants=network-online.target
After=network-online.target

StartLimitInterval=300
StartLimitBurst=3

[Service]
Restart=always
TimeoutSec=30
RestartSec=5
Type=notify
NotifyAccess=main
SyslogIdentifier=%N

User=homepage
Group=homepage

# Polled netdata endpoints — formatted as `<label>=<query-url>[=<click-url>]`,
# joined with commas. Whichever host this is running on gets queried via
# loopback to skip nginx + TLS, with the public FQDN as the click-through
# target; other hosts use the FQDN for both.
{% set _entries = [] %}
{% for h in homepage_netdata_hosts %}
{% set _q = 'http://localhost:19999=' if h == inventory_hostname else '' %}
{% set _ = _entries.append(h ~ '=' ~ _q ~ 'https://netdata.' ~ h ~ '.' ~ domain) %}
{% endfor %}
Environment=NETDATA_HOSTS={{ _entries | join(',') }}
Environment=LISTEN_PORT={{ service_ports.homepage_alerts }}

ExecStart=/usr/bin/python3 /opt/homepage_alerts/alerts_server.py

# Tighten the namespace since this process only needs network + the script
# under /opt. Useful at zero cost — the script is short and well-bounded.
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
