# roles/influxdb/templates/influxdb.service.j2
#jinja2: lstrip_blocks: True
# influxdb.service #######################################################################
[Unit]
Description=influxdb
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=60
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

# No --env TZ: InfluxDB stores timestamps in UTC, so leaving the container in UTC
# keeps log/data correlation trivial (per CLAUDE.md Someday → container timezones).

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name influxdb \
    --init \
    --security-opt=no-new-privileges \
    --log-driver journald \
    --log-opt tag="influxdb" \
    --user {{ influxdb_user.uid }}:{{ influxdb_user.group }}  \
    --env INFLUXD_REPORTING_DISABLED=true \
    --env INFLUXD_STORAGE_VALIDATE_KEYS=true \
    --env INFLUXD_PPROF_DISABLED=true \
    --env INFLUXD_LOG_LEVEL=warn \
    --env DOCKER_INFLUXDB_INIT_MODE=setup \
    --env DOCKER_INFLUXDB_INIT_ORG={{ influxdb_root_org }} \
    --env DOCKER_INFLUXDB_INIT_BUCKET=default \
{% for secret, target, env in [
  ("influxdb_root_user", "influxdb_root_user", "DOCKER_INFLUXDB_INIT_USERNAME_FILE"),
  ("influxdb_root_password", "influxdb_root_password", "DOCKER_INFLUXDB_INIT_PASSWORD_FILE"),
  ("influxdb_admin_token", "influxdb_admin_token", "DOCKER_INFLUXDB_INIT_ADMIN_TOKEN_FILE"),
] %}
{% for arg in secret | podman_secret_file_args(target, env) %}
    {{ arg }} \
{% endfor %}
{% endfor %}
    --volume /mnt/services/influxdb/data:/var/lib/influxdb2 \
    --volume /mnt/services/influxdb/config:/etc/influxdb2 \
    --network homeassistant \
    --publish 127.0.0.1:{{ service_ports.influxdb }}:8086/tcp \
    --health-cmd         {{ "http://localhost:8086/health" | podman_health_curl(head=true) }} \
    --health-startup-cmd {{ "http://localhost:8086/health" | podman_health_curl(head=true) }} \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "30s" \
    --stop-signal SIGINT \
    --stop-timeout 30 \
    {{ influxdb_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
