# roles/systemd_unit/templates/systemd_unit_verify_bad.service.j2
[Unit]
Description=systemd_unit_verify_bad
Bogus=value

[Service]
Type=oneshot
ExecStart=/bin/true
