# roles/systemd_unit/templates/systemd_unit_verify_ok.service.j2
[Unit]
Description=systemd_unit_verify_ok

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecReload=/bin/true

[Install]
WantedBy=multi-user.target
