# roles/udm_dns_failover/templates/on-boot.service
# The /data/on_boot.d/*.sh convention (numeric-prefix ordering, oneshot
# after network-online) is the community unifios-utilities "on-boot-script"
# idiom for surviving UniFi firmware updates. This is a deliberate
# self-contained reimplementation -- don't install the upstream package on
# top, it would collide on this same unit name and directory.
[Unit]
Description=Run /data/on_boot.d scripts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
TimeoutStartSec=120
# set -e so a recovery script that fails fails the unit (visible in
# `systemctl status on-boot`); without it the loop swallows the exit code and
# a broken post-firmware recovery looks like a clean oneshot success. A
# non-executable / glob-no-match entry is in a tested `&&` position, so it
# doesn't trip -e.
ExecStart=/bin/bash -c 'set -euo pipefail; for f in /data/on_boot.d/*.sh; do [ -x "$f" ] && "$f"; done'

[Install]
WantedBy=multi-user.target
