# roles/udm_dns_failover/templates/on-boot.service
[Unit]
Description=Run /data/on_boot.d scripts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
TimeoutStartSec=120
ExecStart=/bin/bash -c 'for f in /data/on_boot.d/*.sh; do [ -x "$f" ] && "$f"; done'

[Install]
WantedBy=multi-user.target
