# roles/udm_dns_failover/templates/udm_keepalived_heartbeat.service.j2
[Unit]
Description=Heartbeat ping for UDM keepalived VRRP peer

[Service]
Type=oneshot
# Ping the success URL when keepalived is up; ping <url>/fail when it is
# down, so a dead-keepalived-but-alive-UDM (the firmware-wiped-/usr case
# this peer exists to catch) raises an immediate, distinct alert instead
# of waiting out the healthchecks.io grace window like an offline UDM.
ExecStart=/bin/bash -c 'if systemctl is-active --quiet keepalived.service; then url="{{ udm_dns_failover_heartbeat_url }}"; else url="{{ udm_dns_failover_heartbeat_url }}/fail"; fi; exec curl -fsS -m 10 --retry 2 --retry-delay 5 --retry-connrefused -o /dev/null "$url"'
