# ansible/roles/pci_passthrough/files/vfio-bind.conf
description "bind your devices to vfio on startup"
version "1.0"

start on starting libvirt-bin
stop on stopped libvirt-bin

script

  # source the config file
  . /etc/vfio-bind.cfg

  /usr/bin/vfio-bind $DEVICES

end script
