# roles/systemd_unit/files/systemd_unit_verify_copy_ok.service
[Unit]
Description=systemd_unit_verify_copy_ok

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
