# roles/zfs_autobackup/files/zfs_autosnapshot.target
[Unit]
Description=Services managed by zfs-autobackup snapshotting

[Install]
WantedBy=default.target
