# roles/netplan/files/netplan_rollback.service
[Unit]
Description=Netplan auto-rollback (fires 90s after a safe-apply if SSH never came back)

[Service]
Type=oneshot
# Keep the oneshot visible as "active (exited)" after a fire so an
# operator (and _verify) can tell the rollback actually ran, instead
# of it reverting to inactive — indistinguishable from "never fired".
RemainAfterExit=yes
SyslogIdentifier=%N
ExecStart=/usr/local/bin/netplan_rollback
