# roles/plex/templates/plex.service.j2
#jinja2: lstrip_blocks: True, trim_blocks: True
# plex.service #######################################################################
[Unit]
Description=plex
After=network-online.target
Wants=network-online.target
# StartLimitInterval must exceed StartLimitBurst × TimeoutSec (3 × 360s =
# 1080s) or the burst limit never fires — three slow starts blow past the
# window and the counter resets. 1800s caps fast-crash loops while still
# letting a large-library DB upgrade ride out a burst's worth of
# TimeoutSec=360 attempts.
StartLimitInterval=1800
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
{% if zfs_has_media_mount %}
After=zfs_mount_mnt_media.service
Requires=zfs_mount_mnt_media.service
{% endif %}
{% if zfs_has_data_mount %}
After=zfs_mount_mnt_data.service
Requires=zfs_mount_mnt_data.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
# Plex migrates its library DB on version bumps; a large library can take
# well past the default ceiling. TimeoutSec must outlast the 300s
# --health-start-period below so the sdnotify=healthy wait wins.
TimeoutSec=360
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name plex \
    --init \
    --security-opt=no-new-privileges \
    --log-driver journald \
    --log-opt tag="plex" \
    --env ADVERTISE_IP="http://{{ external_ips[inventory_hostname] }}:{{ service_ports.plex }}/" \
    --env "TZ={{ system_timezone }}" \
    --env PLEX_UID={{ plex_user.uid }} \
    --env PLEX_GID={{ media_group.gid }} \
    {# The host already owns /config as plex:media, so skip the entrypoint's #}
    {# recursive chown-on-every-start — pure cost (and mtime churn that #}
    {# defeats snapshot dedup) on a large library's Metadata/ tree. #}
    --env CHANGE_CONFIG_DIR_OWNERSHIP=false \
    --volume /mnt/services/plex:/config \
    --volume "/mnt/services/sqlite/plex:/config/Library/Application Support/Plex Media Server/Plug-in Support/Databases" \
    --volume /mnt/media:/media:ro \
{% if zfs_has_data_mount %}
    {# Host-specific extra library roots under /mnt/data (paths exist only #}
    {# on the holding host), declared per-host in plex_extra_library_binds. #}
    {# Gated on the data mount so a host without it never gets podman- #}
    {# auto-created root-owned phantom dirs at these paths. #}
{% for bind in plex_extra_library_binds %}
    --volume "{{ bind.src }}:{{ bind.dest }}:ro" \
{% endfor %}
{% endif %}
    {# 0777, not 0755: Plex runs as the non-root PLEX_UID and the tmpfs #}
    {# mounts root-owned (nothing chowns it), so it must be world-writable #}
    {# for Plex to use it as transcode scratch. tmpfs is container-local, #}
    {# so world-write is benign (jellyfin gets 0755 via its uidmap root). #}
    --mount type=tmpfs,destination=/transcode,tmpfs-size=4G,tmpfs-mode=0777 \
    --network media \
    {# 0.0.0.0, not 127.0.0.1 like jellyfin: native Plex clients (Apple TV, #}
    {# Roku, mobile) and GDM discovery connect straight to :32400 at the #}
    {# ADVERTISE_IP above, bypassing nginx. nginx fronts only the web UI. #}
    --publish 0.0.0.0:{{ service_ports.plex }}:32400/tcp \
{% if has_transcoding_acceleration %}
    --device /dev/dri/renderD128:/dev/dri/renderD128 \
    --group-add={{ getent_group.render[1] | int }} \
{% endif %}
    --health-cmd         '["curl","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:32400/identity"]' \
    --health-startup-cmd '["curl","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:32400/identity"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "300s" \
    --stop-timeout 120 \
    docker.io/plexinc/pms-docker:1.43.2.10687-563d026ea
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
