# roles/paperless/templates/paperless.service.j2
# jinja2: lstrip_blocks: "True"
# paperless.service #######################################################################
[Unit]
Description=paperless
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=240
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
# First level_init pilot. The chain: workload writes "INFO: …" /
# "ERROR: …" / "[259] LOG: …" → level_init scans each line's head,
# prepends the matching <N> syslog priority → podman, with
# --log-driver=none, keeps container stdio attached to its own
# (since we dropped --detach) → systemd reads podman's stdout/stderr,
# stream-parses <N>, sets PRIORITY before the line lands in the
# journal. SyslogIdentifier=%N attributes the entry to this unit.
#
# Driver choice: `none` (not `journald` or `k8s-file`):
#  - journald has podman/conmon double-log every line, once raw via
#    sd_journal_send (hardcoded PRIORITY-by-stream) and once via the
#    systemd stream parser. Duplicates with mismatched priorities.
#  - k8s-file works correctly via the same systemd path but also
#    writes a per-container log file under /var/lib/containers/...
#    that has to be size-bounded.
#  - none does what we want: nothing duplicates, nothing on disk,
#    journal is the single source of truth.
#
# Cost of `none`: `podman logs paperless` returns nothing. Acceptable
# since journalctl -u paperless shows the same data with corrected
# priorities.
#
# --detach is dropped because passthrough/none require podman to stay
# foreground so its stdio remains connected to systemd's pipe. With
# --detach, podman exits and the bytes vanish.
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --replace \
    --rm \
    --name paperless \
    --init \
    --init-path=/usr/local/libexec/level_init \
    --log-driver=none \
    --volume /mnt/services/paperless/data:/usr/src/paperless/data \
    --volume /mnt/services/sqlite/paperless:/mnt/services/sqlite/paperless \
    --volume /mnt/services/paperless/media:/usr/src/paperless/media \
    --volume /mnt/services/paperless/consume:/usr/src/paperless/consume \
    --volume /mnt/services/paperless/export:/usr/src/paperless/export \
    --network redis \
    --publish 127.0.0.1:8001:8000/tcp \
    --env USERMAP_UID={{ paperless_user.uid }} \
    --env USERMAP_GID={{ paperless_user.group }} \
    --env PAPERLESS_REDIS="redis://redis:6379/1" \
    --env PAPERLESS_OCR_LANGUAGES="eng fra nld" \
    --env PAPERLESS_OCR_LANGUAGE=fra \
    --env PAPERLESS_URL=https://paperless.{{ inventory_hostname }}.{{ domain }} \
    --secret=paperless_secret_key,type=env,target=PAPERLESS_SECRET_KEY \
    --health-cmd         "curl --head --location --fail --silent --show-error --connect-timeout 1 --max-time 5 -o /dev/null http://localhost:8000/" \
    --health-startup-cmd "curl --head --location --fail --silent --show-error --connect-timeout 1 --max-time 5 -o /dev/null http://localhost:8000/" \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 120 \
    ghcr.io/paperless-ngx/paperless-ngx:2.12
# ExecStop / ExecStopPost are unnecessary in this foreground model:
# systemd's SIGTERM to MainPID (podman) is sig-proxied to the
# container's PID 1, the workload exits gracefully, podman exits,
# ExecStart returns. Whatever ExecStop=podman-stop would do is
# already done by the time it'd run. --replace --rm on ExecStart
# handles the rare case of a leftover container after a hard kill.

[Install]
WantedBy=default.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
