# roles/paperless/templates/paperless.service.j2
#jinja2: lstrip_blocks: True
# paperless.service #######################################################################
[Unit]
Description=paperless
After=network-online.target
Wants=network-online.target
# redis.service is Type=notify --sdnotify=healthy, so it only reaches
# `active` once its in-container redis-cli ping passes. Ordering after it
# (Wants, not Requires — paperless tolerates redis bouncing post-start via
# its own retry loop) means paperless boots against a healthy broker
# instead of spamming connection retries until redis comes up. Also
# guarantees the shared `redis` podman network exists before paperless
# attaches to it.
After=redis.service
Wants=redis.service
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=240
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name paperless \
    --log-driver journald \
    --log-opt tag="paperless" \
    --security-opt=no-new-privileges \
    --volume /mnt/services/paperless/data:/usr/src/paperless/data \
    --volume /mnt/services/sqlite/paperless:/mnt/services/sqlite/paperless \
    --volume /mnt/services/paperless/media:/usr/src/paperless/media \
    --volume /mnt/services/paperless/consume:/usr/src/paperless/consume \
    --volume /mnt/services/paperless/export:/usr/src/paperless/export \
    --network redis \
    --publish 127.0.0.1:{{ service_ports.paperless }}:8000/tcp \
    {# USERMAP_UID/GID (CLAUDE.md "User namespacing" tier 2): paperless-ngx's
       s6 init usermods its baked-in paperless user to these ids and drops to
       it. Tier 2 fits a droppable image -- the fake-root --uidmap pattern is
       reserved for images that insist on id 0 (hyperdx). #}
    --env USERMAP_UID={{ paperless_user.uid }} \
    --env USERMAP_GID={{ paperless_user.group }} \
    --env PAPERLESS_REDIS="redis://:{{ redis_password }}@redis.dns.podman:6379/1" \
    --env PAPERLESS_DB_TIMEOUT=30 \
    --env PAPERLESS_OCR_LANGUAGES="eng fra nld" \
    --env PAPERLESS_OCR_LANGUAGE=fra \
    --env PAPERLESS_TIME_ZONE={{ system_timezone }} \
    --env PAPERLESS_URL=https://paperless.{{ inventory_hostname }}.{{ domain }} \
    --env PAPERLESS_APPS=allauth.socialaccount.providers.openid_connect \
    --env PAPERLESS_ACCOUNT_DEFAULT_HTTP_PROTOCOL=https \
    --env PAPERLESS_ACCOUNT_EMAIL_VERIFICATION=none \
    --env PAPERLESS_SOCIAL_AUTO_SIGNUP=true \
    --env PAPERLESS_SOCIAL_ACCOUNT_SYNC_GROUPS=true \
    --env PAPERLESS_REDIRECT_LOGIN_TO_SSO=true \
    --env PAPERLESS_DISABLE_REGULAR_LOGIN=true \
    --env PAPERLESS_SOCIALACCOUNT_PROVIDERS_FILE=/run/secrets/socialaccount_providers.json \
    --env PAPERLESS_SECRET_KEY_FILE=/run/secrets/secret_key \
    --secret=paperless_secret_key,type=mount,target=secret_key \
    --secret=paperless_socialaccount_providers,type=mount,target=socialaccount_providers.json \
    --health-cmd         '["curl","--head","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:8000/"]' \
    --health-startup-cmd '["curl","--head","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:8000/"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 120 \
    ghcr.io/paperless-ngx/paperless-ngx:2.20.15
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
