# roles/udm_dns_failover/files/udm_keepalived_heartbeat.timer
[Unit]
Description=Heartbeat timer for UDM keepalived VRRP peer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
