# roles/podman/files/podman-prune.service
# vim:ft=systemd-unit-file
[Unit]
Description=Prune unused podman containers, images, networks and volumes

[Service]
Type=oneshot
# No -a: keeps named (tagged) images warm so the next pull is a no-op.
# --volumes drops anonymous volumes not referenced by any container.
# until=48h spares anything created in the last two days so a recent
# stopped container or fresh dangling layer survives one daily pass.
ExecStart=/usr/bin/podman system prune -f --volumes --filter until=48h
