# roles/nexus/templates/nexus.service.j2
#jinja2: lstrip_blocks: True
# nexus.service #######################################################################
[Unit]
Description=nexus
After=network-online.target
Wants=network-online.target
# StartLimitIntervalSec must comfortably exceed StartLimitBurst × TimeoutStartSec
# (here 3 × 240s = 720s) or the burst limit never fires -- three slow starts
# blow past the window and the counter resets, so a wedged JVM restarts forever
# instead of landing in `failed` where netdata's systemdunits collector sees it.
StartLimitIntervalSec=1800
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
{% endif %}
{% if zfs_has_scratch_mount %}
After=zfs_mount_mnt_scratch.service
Requires=zfs_mount_mnt_scratch.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutStartSec=240
# TimeoutStopSec exceeds --stop-timeout 120 below so a slow graceful container
# stop doesn't eat into the start budget.
TimeoutStopSec=150
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

# Nexus's DB + config ride the snapshotted services dataset (/nexus-data).
# Only the `default` blob store is bind-mounted to /mnt/scratch (NOT
# snapshotted): it holds the reconstructible proxy caches, so keeping them
# out of snapshots keeps replication small. Any *named* blob store created
# under /nexus-data/blobs instead lands on the snapshotted services dataset
# -- point hosted/irreplaceable repos (their artifacts can't be re-proxied)
# at one of those. After a scratch loss the default cache is cold and
# re-proxies on demand; a DB restore against an empty default store is
# expected, not corruption.
ExecStartPre=/bin/rm -f %t/%n.ctr-id
# INSTALL4J_ADD_VM_PARAMS relocates the java.util.prefs userRoot onto the
# /user-root volume so it sits beside encryption.json, off the snapshotted
# /nexus-data mount. The healthcheck uses the Alpine base image's busybox
# wget (a -ubi image tag would not ship it).
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name nexus \
    --init \
    --security-opt=no-new-privileges \
    --network nexus \
    --log-driver journald \
    --log-opt tag="nexus" \
    --env INSTALL4J_ADD_VM_PARAMS="{{ nexus_jvm_params }} -Djava.util.prefs.userRoot=/user-root" \
    --env NEXUS_SECRETS_KEY_FILE=/user-root/encryption.json \
    --user 200:200 \
    {{ nexus_user | podman_idmap_args(container_uid=200) | join(" \\\n    ") }} \
    --volume /mnt/services/nexus/data:/nexus-data \
    --volume /mnt/scratch/nexus:/nexus-data/blobs/default \
    --volume /mnt/services/nexus/user:/user-root \
    --publish 127.0.0.1:{{ service_ports.nexus }}:8081/tcp \
    --health-cmd         {{ "http://localhost:8081/service/rest/v1/status/writable" | podman_health_wget }} \
    --health-startup-cmd {{ "http://localhost:8081/service/rest/v1/status/writable" | podman_health_wget }} \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "180s" \
    --stop-timeout 120 \
    {{ nexus_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
