# ansible/roles/services/files/systemd/couchpotato.service
# couchpotato.service #######################################################################
[Unit]
Description=Couchpotato
After=docker.service mnt-media.mount zfs-mount.service sabnzbd.service
Requires=docker.service mnt-media.mount zfs-mount.service sabnzbd.service
PartOf=compose.service

StartLimitInterval=300
StartLimitBurst=5

[Service]
Restart=always
TimeoutSec=300
RestartSec=15

ExecStartPre=-/usr/bin/docker kill couchpotato
ExecStartPre=-/usr/bin/docker rm couchpotato
ExecStartPre=/usr/bin/docker pull akohlbecker/couchpotato
ExecStartPre=/usr/bin/docker create \
    --init \
    --env-file /etc/env.d/couchpotato \
    --name couchpotato \
    --volume /mnt/docker/couchpotato:/data \
    --volume /mnt/media:/media \
    --label "traefik.docker.network=traefik" \
    --label "traefik.enable=true" \
    --label "traefik.frontend.rule=Host:couchpotato.kohlby.fr" \
    --label "traefik.port=5050" \
    --label "traefik.protocol=http" \
    --label "traefik.frontend.auth.forward.address=http://traefik-auth-cloudflare/auth/5e02944c3d8ef0cc0d9f2eb0906557b75ecf47ac3bdcd94aa964cd5431d35996" \
    akohlbecker/couchpotato
ExecStartPre=/usr/bin/docker network connect usenet couchpotato
ExecStartPre=/usr/bin/docker network connect traefik couchpotato
ExecStart=/usr/bin/docker start --attach couchpotato
ExecStop=/usr/bin/docker stop couchpotato

[Install]
WantedBy=compose.service
