# ansible/roles/services/files/systemd/docker-hub-schedule.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-hub-schedule
