# ansible/roles/services/files/systemd/docker-hub-schedule.timer
[Timer]
OnUnitActiveSec=1w
Persistent=true

[Install]
WantedBy=compose.service
