# ansible/roles/services/files/systemd/usenet-maintenance.path
[Path]
PathModified=/mnt/services/headphones/config.ini
PathModified=/mnt/services/sabnzbd/sabnzbd.ini
PathModified=/mnt/services/couchpotato/settings.conf

[Install]
WantedBy=compose.service
