# ansible/roles/services/files/systemd/usenet-maintenance.timer
[Timer]
OnBootSec=15min
OnUnitActiveSec=24h

[Install]
WantedBy=compose.service
