# ansible/roles/services/files/systemd/usenet-maintenance.path
[Path]
PathModified=/mnt/docker/sickrage/config.ini
PathModified=/mnt/docker/sabnzbd/sabnzbd.ini
PathModified=/mnt/docker/couchpotato/settings.conf

[Install]
WantedBy=compose.service
