# roles/github_runner/templates/github_runner.service.j2
# jinja2: lstrip_blocks: "True"
# github_runner@.service ###############################################
# Templated unit; one instance per entry in github_runner_instances on
# this host. Instance name = <repo>_<suffix> (e.g. homelab_vm1). Per-
# instance config (REPO_URL, RUNNER_NAME, RUNNER_LABELS) ships via
# /etc/default/github_runner@%i; first-start registration token ships
# via the github_runner_token_%i podman secret. State persists on the
# /mnt/services/github_runner/%i bind-mount so container restarts and
# image bumps don't re-register.
[Unit]
Description=github-actions runner (%i)
After=network-online.target
Wants=network-online.target
# actions/runner inside the container spawns sibling workflow containers
# via the bind-mounted rootless podman socket; wait for github_runner's
# user manager (which owns that socket) to be up before starting.
After=user@{{ github_runner_user.uid }}.service
Requires=user@{{ github_runner_user.uid }}.service
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_scratch_mount %}
After=zfs_mount_mnt_scratch.service
Requires=zfs_mount_mnt_scratch.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
RestartSec=5
Delegate=yes
# Cold-start has to fit: image pull (first time), container create, the
# entrypoint's curl POST to GitHub's generate-jitconfig endpoint, and
# Runner.Listener's startup to "Listening for jobs" (the healthcheck's
# ready signal). 300s covers the slowest of these (network-bound pull)
# with headroom.
TimeoutStartSec=300
TimeoutStopSec=60
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

User=github_runner
Group=github_runner

# %t resolves to /run on a system unit (system runtime dir), regardless
# of User=. podman's --cidfile would then try to write directly under
# /run, which github_runner can't. RuntimeDirectory=%N tells systemd
# to mkdir /run/<unit-name> with mode 0755 owned by User=/Group=
# (here github_runner:github_runner), and tear it down when the unit
# stops. The --cidfile / --cidfile-cleanup paths below all reference
# %t/%N/ctr-id to land inside it. Sibling podman roles (healthchecks,
# openobserve, ...) skip this dance because they run the unit as root
# at the systemd layer; we keep User= so a podman crash + escape
# lands at github_runner, not root.
RuntimeDirectory=%N

# Wait for the github_runner user manager to bring its rootless
# podman.socket up. Bounded 60 x 1s, well inside TimeoutStartSec. If the
# socket never appears, the journal records the explicit cause rather
# than a generic start-timeout.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -S /run/user/{{ github_runner_user.uid }}/podman/podman.sock ] && exit 0; sleep 1; done; echo "podman.sock not present after 60s" >&2; exit 1'

# --user / --uidmap / --gidmap: fake-root pattern (see CLAUDE.md
#   "Container User Namespacing"). In-container uid 0 maps to host
#   github_runner via the +0 override; everything else falls through
#   1:1 from the base 65536-range mapping. Writes to the bound
#   podman.sock and the /work bind-mount land owned by host
#   github_runner. NOT `--userns keep-id`: that would map the *current
#   rootless user* to the same in-container uid, leaving in-container
#   root pointing at an unmapped subuid and failing every write.
# --volume podman.sock: DooD -- actions/runner's docker-CLI shell-outs
#   (for workflows with `container:` blocks) hit the host's rootless
#   podman, NOT a daemon inside the runner container.
# --volume /work: actions/runner's checkout root (_work tree). Lives on
#   /mnt/scratch (terabytes) so per-job checkouts don't pressure the
#   container's writable layer. The actions-runner binaries themselves
#   live in /actions-runner inside the image -- ephemeral runners run
#   straight out of the image with no host-side state volume.
# --env-file: per-instance REPO_OWNER / REPO_NAME / RUNNER_NAME /
#   RUNNER_LABELS / GITHUB_PAT_FILE, written by the role to
#   /etc/default/github_runner@%i.
# --secret github_runner_pat: long-lived fine-grained PAT consumed by
#   the entrypoint on every container start (mints a fresh JIT config
#   via /actions/runners/generate-jitconfig). Mount-type, not env-type,
#   so the PAT never lands in /proc/<pid>/environ visible to any
#   sibling workflow container that has the DooD socket. Same PAT
#   shared across every instance on the host (API call is per-repo).
# --sdnotify=healthy: unit's `active` state gates on the healthcheck,
#   not on conmon's exec. The container is ready iff Runner.Listener
#   is running and polling -- exactly what systemctl status should
#   reflect, and what removes the need for an ansible-side wait_for.
# --health-cmd / startup-cmd: actions/runner exposes no HTTP probe, so
#   we go with CLAUDE.md "Container Healthchecks" option 3 (service-
#   native CLI): pgrep the Runner.Listener process. procps ships in
#   ubuntu:24.04's base so the runner image inherits it without an
#   explicit apt install. When the listener loses its GitHub
#   connection it tries to reconnect; if reconnection fails it exits
#   non-zero, the container exits, and systemd's Restart=always
#   respawns -- so the healthcheck is the "is it alive" signal
#   surfaced in `systemctl status`/downstream monitoring, while the
#   exit+respawn loop drives the actual recovery.
# --health-start-period 30s: first-start runs the entrypoint's curl
#   POST to GitHub before run.sh. 30s leaves headroom for a slow API
#   response without flapping the unit unhealthy on boot.
ExecStartPre=/bin/rm -f %t/%N/ctr-id
ExecStart=/usr/bin/podman run \
  --cidfile=%t/%N/ctr-id \
  --cgroups=split \
  --sdnotify=healthy \
  --detach \
  --replace \
  --rm \
  --name github_runner_%i \
  --log-driver journald \
  --log-opt tag="github_runner_%i" \
  --user 0:0 \
  --uidmap=0:0:65536 \
  --uidmap=+0:{{ github_runner_user.uid }}:1 \
  --gidmap=0:0:65536 \
  --gidmap=+0:{{ github_runner_user.group }}:1 \
  --volume /run/user/{{ github_runner_user.uid }}/podman/podman.sock:/var/run/docker.sock \
  --volume /mnt/scratch/github_runner/workdir_%i:/work:rw \
  --env-file /etc/default/github_runner@%i \
  --secret github_runner_pat,type=mount,target=/run/secrets/github_pat \
  --health-cmd         "pgrep -f Runner.Listener >/dev/null" \
  --health-startup-cmd "pgrep -f Runner.Listener >/dev/null" \
  --health-startup-interval "1s" \
  --health-startup-timeout "5s" \
  --health-start-period "30s" \
  --stop-timeout 30 \
  nexus.lab.fahm.fr/homelab/runner:latest
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%N/ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%N/ctr-id

[Install]
WantedBy=default.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
