# roles/github_runner/templates/github_runner.service.j2
#jinja2: lstrip_blocks: True
# github_runner@.service ###############################################
# Templated system-scope unit running as root, demoting to host
# github_runner inside the container via the +0 uidmap override.
# One instance per entry in github_runner_instances on this host.
# Instance name = <repo>_<suffix> (e.g. homelab_vm1). Per-instance
# config ships via /etc/default/github_runner@%i; the long-lived
# GitHub PAT ships via the github_runner_pat podman secret (host-
# level, shared across instances). Ephemeral + JIT: the entrypoint
# mints a single-use JIT config on every container start and run.sh
# self-deregisters when the job completes.
#
# Same shape as every other system-scope podman role on this host
# (healthchecks, openobserve, etc.): no User= directive, podman runs
# as root, and the in-container processes are demoted to host
# github_runner via --uidmap=+0:<uid>:1. Two earlier variants in
# this branch (system-scope-with-User=github_runner, then user-scope)
# tried to keep podman itself rootless for defense in depth; both
# hit failures (cgroup permissions for exec, then chmod EPERM from
# the rootless-on-rootless userns layering against the bind-mounted
# /work). Going with the standard pattern accepts a slightly wider
# blast radius for a kernel-vuln-class container escape (lands as
# host root vs host github_runner) in exchange for a simpler,
# better-tested code path. Acceptable for homelab CI per CLAUDE.md.
#
# Sibling workflow containers (spawned via the DooD socket from
# inside the runner) STAY rootless under github_runner -- the DooD
# bind mount points at /run/user/<uid>/podman/podman.sock, which is
# served by github_runner's user-scope podman.socket. So the
# defense-in-depth story holds for the second-level containers,
# which is where untrusted action code actually runs.
[Unit]
Description=github-actions runner (%i)
# DooD socket lives at /run/user/<uid>/podman/podman.sock under
# github_runner's user manager. A system-scope unit can't directly
# After= a user-scope socket, but it can After= the user manager
# itself (user@<uid>.service is a system unit). linger keeps the
# user manager up at boot and podman.socket is enabled in user
# scope, so by the time user@<uid>.service is `active`, podman.sock
# is reachable.
After=user@{{ github_runner_user.uid }}.service
Requires=user@{{ github_runner_user.uid }}.service
# Runner image lives at nexus.lab.fahm.fr/homelab/runner:latest. nexus is
# contacted on cold-start (--pull=missing fires once on a fresh host) and
# from github_runner_image_rotate.timer's 5-minute tick; steady-state
# restarts after a job land on the cached layer and never touch nexus.
# After= orders cold-start so the first pull doesn't race nexus bringup
# at boot; Wants= (not Requires=) keeps a nexus stop from cascading into
# in-flight runner jobs. On hosts without nexus.service the After= is
# silently a no-op, so the unit stays portable.
After=nexus.service
Wants=nexus.service
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_scratch_mount %}
After=zfs_mount_mnt_scratch.service
Requires=zfs_mount_mnt_scratch.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
# RuntimeDirectory= creates /run/%N (= /run/github_runner@<inst>.service/)
# on start and removes it on stop -- cidfile lands there. %t is /run
# for system units so the cidfile path is %t/%N/%n.ctr-id.
RuntimeDirectory=%N
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
RestartSec=5
Delegate=yes
# Cold-start has to fit: image pull (first time), container create, the
# entrypoint's curl POST to GitHub's generate-jitconfig endpoint, and
# Runner.Listener's startup to "Listening for jobs" (the healthcheck's
# ready signal). 300s covers the slowest of these (network-bound pull)
# with headroom.
TimeoutStartSec=300
# TimeoutStopSec has to exceed --stop-timeout below by enough margin
# for podman to SIGKILL the container, clean up the storage graph
# entry, and exit. 30s headroom over the 60s --stop-timeout (chosen
# below) covers the slowest observed cleanup.
TimeoutStopSec=90
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

# --cgroups=split: container + conmon live as sub-cgroups under the
#   systemd unit's cgroup tree. Resource limits / accounting in the
#   unit apply to the container; systemctl status shows real resource
#   usage; systemctl kill reaps the container. Standard system-scope
#   + Delegate=yes idiom, same as every other podman service role.
# --init: PID 1 inside the container would otherwise be the entrypoint
#   shell, which doesn't reap orphaned children. actions/runner spawns
#   Runner.Listener -> Runner.Worker -> per-step processes (node, git,
#   ...); a workflow step that double-forks and exits leaves zombies
#   the shell-as-PID-1 won't collect, slowly burning kernel PID slots
#   over the container's lifetime. --init wires podman's catatonit at
#   PID 1 to handle signal forwarding and zombie reaping. Compatible
#   with --sdnotify=healthy (podman owns the NOTIFY_SOCKET path, not
#   PID 1).
# --user / --uidmap / --gidmap: fake-root pattern from CLAUDE.md
#   "Container User Namespacing". podman runs as real root (no
#   User= in the unit), so the base 0:0:65536 would map in-container
#   uid 0 to host root. The +0:<github_runner.uid>:1 override
#   demotes in-container uid 0 to host github_runner so writes to
#   /work and the DooD socket land at the right owner without
#   granting real-root inside the container.
# --volume podman.sock: DooD -- actions/runner's docker-CLI shell-outs
#   (for workflows with `container:` blocks) hit github_runner's
#   rootless podman, NOT a daemon inside the runner container.
# --volume /opt/actions-runner: path-mirrored from the host (where
#   roles/github_runner downloaded + extracted the actions/runner
#   tarball). The runner binaries live here on both sides at the same
#   path, so when actions/runner spawns docker create with
#   -v /opt/actions-runner/externals:/__e:ro for a workflow container,
#   docker (resolving the source on the host via the bind-mounted
#   socket) finds real files instead of EACCES on a missing path.
#   rw (not ro) because actions/runner writes _diag/ logs into the
#   install dir at runtime.
# --volume /mnt/.../workdir_%i: path-mirrored at the same path on
#   both sides for the same DooD-bind-mount reason. actions/runner's
#   _work tree lives here, plus _temp/_actions/_tool subdirs that
#   workflow containers see at /__w/_* via DooD bind mounts that
#   reference this host path. Lives on /mnt/scratch (terabytes) so
#   per-job checkouts don't pressure the container's writable layer.
#   RUNNER_WORK_FOLDER env below tells actions/runner to use this
#   path as its work root (vs the default _work next to run.sh).
# --env-file: per-instance REPO_OWNER / REPO_NAME / RUNNER_NAME /
#   RUNNER_LABELS / GITHUB_PAT_FILE, written by the role to
#   /etc/default/github_runner@%i.
# --env CONTAINER_HOST / DOCKER_HOST: point in-container podman /
#   docker CLIs at the bind-mounted DooD socket. Without these, a
#   workflow step that runs `podman build` (or `docker build`) tries
#   to do nested rootless podman inside the runner container -- which
#   fails because the container doesn't have CAP_SYS_ADMIN to do its
#   own mount-namespace setup (errors with "remount /, flags: 0x44000:
#   permission denied" during ApplyLayer). With these set, the CLI
#   transparently uploads the build context to the host's rootless
#   podman over the socket, builds happen there (where mount-ns
#   operations have the privileges they need), and resulting images
#   land in github_runner's host-side storage. Setting both var names
#   covers podman CLI (uses CONTAINER_HOST) and the podman-docker
#   alias / any real docker client (uses DOCKER_HOST).
# --secret github_runner_pat: long-lived fine-grained PAT consumed by
#   the entrypoint on every container start (mints a fresh JIT config
#   via /actions/runners/generate-jitconfig). Mount-type, not env-type,
#   so the PAT never lands in /proc/<pid>/environ visible to any
#   sibling workflow container that has the DooD socket. Same PAT
#   shared across every instance on the host (API call is per-repo).
# --sdnotify=healthy: unit's `active` state gates on the healthcheck,
#   not on conmon's exec. The container is ready iff Runner.Listener
#   is running and polling -- exactly what systemctl status should
#   reflect, and what removes the need for an ansible-side wait_for.
# --health-cmd / startup-cmd: actions/runner exposes no HTTP probe, so
#   we go with CLAUDE.md "Container Healthchecks" option 3 (service-
#   native CLI): pgrep the Runner.Listener process. procps ships in
#   ubuntu:24.04's base so the runner image inherits it without an
#   explicit apt install. When the listener loses its GitHub
#   connection it tries to reconnect; if reconnection fails it exits
#   non-zero, the container exits, and systemd's Restart=always
#   respawns -- so the healthcheck is the "is it alive" signal
#   surfaced in `systemctl status`/downstream monitoring, while the
#   exit+respawn loop drives the actual recovery. Using -x (exact
#   match against /proc/<pid>/comm) over -f (full cmdline substring)
#   so a workflow process whose argv happens to mention "Runner.Listener"
#   can't keep the unit reporting healthy. "Runner.Listener" is 15
#   chars, within TASK_COMM_LEN (16).
# --health-start-period 30s: first-start runs the entrypoint's curl
#   POST to GitHub before run.sh. 30s leaves headroom for a slow API
#   response without flapping the unit unhealthy on boot.
# --stop-timeout 60: actions/runner's main loop is a long-poll against
#   GitHub's GET /runners/messages (~50s server hold). SIGTERM only
#   takes effect once the in-flight long-poll completes, so a smaller
#   timeout fires SIGKILL prematurely and noisy "StopSignal SIGTERM
#   failed to stop container" warnings flood the journal on every
#   restart cascade (template edit -> all instances bounce). 60s
#   covers the long-poll plus a little graceful-exit wiggle. Kept
#   strictly less than TimeoutStopSec above so podman wins the race
#   and systemd doesn't SIGKILL podman mid-cleanup.
# --pull defaults to 'missing': only pull when the image isn't in the
#   local store (i.e. first start on a fresh host). Subsequent rotation
#   onto a newer runner-image build is driven by
#   github_runner_image_rotate.timer, which polls nexus every 5 minutes,
#   refreshes the local store via `podman pull`, and restarts idle stale
#   instances. Putting --pull=newer back here would just duplicate the
#   manifest check on every container start (i.e. after every job
#   completion) for no benefit -- the watcher has already cached the
#   newest image locally, so the cached layer satisfies the next start
#   either way.
# --security-opt seccomp=unconfined: podman 4.9.3 (ubuntu jammy host)
#   ships a default seccomp profile that predates the kernel 6.6
#   fchmodat2 syscall (number 452). Unknown syscalls return EPERM by
#   default, which breaks GNU tar 1.35+ inside the container: tar uses
#   fchmodat2 with AT_SYMLINK_NOFOLLOW for safe directory mode
#   restoration during extract, and every actions/checkout download
#   trips it -- "Cannot change mode to rwxrwxr-x: Operation not
#   permitted" on every dir, non-zero exit, every job fails. Unconfining
#   seccomp restores the chmod path; the threat model is unchanged
#   because the +0 uidmap demotion already caps in-container compromise
#   at unprivileged host github_runner. A narrower fix would be a
#   custom seccomp.json that just adds fchmodat2 to the allowlist;
#   acceptable for homelab CI to skip that and unconfine wholesale.
ExecStartPre=/bin/rm -f %t/%N/%n.ctr-id
ExecStart=/usr/bin/podman run \
  --cidfile=%t/%N/%n.ctr-id \
  --cgroups=split \
  --security-opt seccomp=unconfined \
  --sdnotify=healthy \
  --detach \
  --replace \
  --rm \
  --init \
  --name github_runner_%i \
  --init \
  --log-driver journald \
  --log-opt tag="github_runner_%i" \
  --user 0:0 \
  --uidmap=0:0:65536 \
  --uidmap=+0:{{ github_runner_user.uid }}:1 \
  --gidmap=0:0:65536 \
  --gidmap=+0:{{ github_runner_user.group }}:1 \
  --volume %t/user/{{ github_runner_user.uid }}/podman/podman.sock:/var/run/docker.sock \
  --volume /opt/actions-runner:/opt/actions-runner:rw \
  --volume /mnt/scratch/github_runner/workdir_%i:/mnt/scratch/github_runner/workdir_%i:rw \
  --env-file /etc/default/github_runner@%i \
  --env CONTAINER_HOST=unix:///var/run/docker.sock \
  --env DOCKER_HOST=unix:///var/run/docker.sock \
  --env RUNNER_WORK_FOLDER=/mnt/scratch/github_runner/workdir_%i \
  --secret github_runner_pat,type=mount,target=/run/secrets/github_pat \
  --health-cmd         '["pgrep","--count","--exact","Runner.Listener"]' \
  --health-startup-cmd '["pgrep","--count","--exact","Runner.Listener"]' \
  --health-startup-interval "1s" \
  --health-startup-timeout "5s" \
  --health-start-period "30s" \
  --stop-timeout 60 \
  nexus.lab.fahm.fr/homelab/runner:latest
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%N/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%N/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
