# roles/github_runner/templates/github_runner.service.j2
# jinja2: lstrip_blocks: "True"
# github_runner@.service ###############################################
# Templated user-scope unit; one instance per entry in
# github_runner_instances on this host, all running under
# github_runner's user@<uid>.service. Instance name = <repo>_<suffix>
# (e.g. homelab_vm1). Per-instance config ships via
# /etc/default/github_runner@%i; the long-lived GitHub PAT ships via
# the github_runner_pat podman secret (host-level, shared across
# instances). Ephemeral + JIT: the entrypoint mints a single-use JIT
# config on every container start and run.sh self-deregisters when
# the job completes.
#
# User scope (not system scope, vs other podman roles like
# healthchecks/openobserve): rootless podman + custom uidmaps need
# cgroup writes within a user-owned tree for `podman exec` and
# `podman healthcheck` to work. A system-scope unit running as
# User=github_runner gets a Delegate=yes tree under system.slice/
# but the split sub-cgroups end up not-quite-writable enough for
# exec. Under github_runner's user manager the entire tree lives
# under user.slice/user-<uid>.slice/ -- fully owned by
# github_runner, no Delegate dance, exec joins cleanly.
[Unit]
Description=github-actions runner (%i)
# Order after the github_runner user manager's podman.socket so the
# bind-mounted DooD socket below resolves at start time. Requires=
# pulls podman.socket up if it's not already; After= orders us
# after its readiness. Both are user-scope (no --user qualifier
# needed inside a user-scope unit -- the manager already is the
# user's manager).
After=podman.socket
Requires=podman.socket
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_scratch_mount %}
# zfs deps are cross-scope: this is a user unit, zfs_mount_*.service
# / zfs_autosnapshot.target are system units, and systemd silently
# ignores cross-scope After=/Requires=/PartOf= edges. The lines
# below are kept as intent-documentation -- the actual ordering
# (don't start until /mnt/scratch is mounted; stop+restart around
# zfs autosnapshot) is delegated to a drop-in on user@<uid>.service
# in main.yml so the whole user manager (and all its child units,
# including every github_runner@<inst>.service) waits transitively.
After=zfs_mount_mnt_scratch.service
Requires=zfs_mount_mnt_scratch.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
RestartSec=5
Delegate=yes
# Cold-start has to fit: image pull (first time), container create, the
# entrypoint's curl POST to GitHub's generate-jitconfig endpoint, and
# Runner.Listener's startup to "Listening for jobs" (the healthcheck's
# ready signal). 300s covers the slowest of these (network-bound pull)
# with headroom.
TimeoutStartSec=300
# TimeoutStopSec has to exceed --stop-timeout below by enough margin
# for podman to SIGKILL the container, clean up the rootless storage
# graph entry, and exit. 30s headroom over the 60s --stop-timeout
# (chosen below) covers the slowest observed cleanup.
TimeoutStopSec=90
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

# %t in a user-scope unit resolves to /run/user/<uid>/, which is
# already github_runner-owned -- so --cidfile=%t/%n.ctr-id lands
# in a writable path with no RuntimeDirectory trick needed.

# --user / --uidmap / --gidmap: fake-root pattern (see CLAUDE.md
#   "Container User Namespacing"). In-container uid 0 maps to host
#   github_runner via the +0 override; everything else falls through
#   1:1 from the base 65536-range mapping. Writes to the bound
#   podman.sock and the /work bind-mount land owned by host
#   github_runner. NOT `--userns keep-id`: that would map the *current
#   rootless user* to the same in-container uid, leaving in-container
#   root pointing at an unmapped subuid and failing every write.
# --cgroups=split: container + conmon live as sub-cgroups under the
#   systemd unit's cgroup tree. Resource limits / accounting in the
#   unit apply to the container; systemctl status shows real
#   resource usage; systemctl kill reaps the container. Works here
#   because the unit is user-scope and the cgroup tree lives entirely
#   under user.slice/user-<uid>.slice/, all github_runner-owned.
# --volume podman.sock: DooD -- actions/runner's docker-CLI shell-outs
#   (for workflows with `container:` blocks) hit the host's rootless
#   podman, NOT a daemon inside the runner container.
# --volume /work: actions/runner's checkout root (_work tree). Lives on
#   /mnt/scratch (terabytes) so per-job checkouts don't pressure the
#   container's writable layer. The actions-runner binaries themselves
#   live in /actions-runner inside the image -- ephemeral runners run
#   straight out of the image with no host-side state volume.
# --env-file: per-instance REPO_OWNER / REPO_NAME / RUNNER_NAME /
#   RUNNER_LABELS / GITHUB_PAT_FILE, written by the role to
#   /etc/default/github_runner@%i.
# --secret github_runner_pat: long-lived fine-grained PAT consumed by
#   the entrypoint on every container start (mints a fresh JIT config
#   via /actions/runners/generate-jitconfig). Mount-type, not env-type,
#   so the PAT never lands in /proc/<pid>/environ visible to any
#   sibling workflow container that has the DooD socket. Same PAT
#   shared across every instance on the host (API call is per-repo).
# --sdnotify=healthy: unit's `active` state gates on the healthcheck,
#   not on conmon's exec. The container is ready iff Runner.Listener
#   is running and polling -- exactly what systemctl status should
#   reflect, and what removes the need for an ansible-side wait_for.
# --health-cmd / startup-cmd: actions/runner exposes no HTTP probe, so
#   we go with CLAUDE.md "Container Healthchecks" option 3 (service-
#   native CLI): pgrep the Runner.Listener process. procps ships in
#   ubuntu:24.04's base so the runner image inherits it without an
#   explicit apt install. When the listener loses its GitHub
#   connection it tries to reconnect; if reconnection fails it exits
#   non-zero, the container exits, and systemd's Restart=always
#   respawns -- so the healthcheck is the "is it alive" signal
#   surfaced in `systemctl status`/downstream monitoring, while the
#   exit+respawn loop drives the actual recovery.
# --health-start-period 30s: first-start runs the entrypoint's curl
#   POST to GitHub before run.sh. 30s leaves headroom for a slow API
#   response without flapping the unit unhealthy on boot.
# --stop-timeout 60: actions/runner's main loop is a long-poll against
#   GitHub's GET /runners/messages (~50s server hold). SIGTERM only
#   takes effect once the in-flight long-poll completes, so a smaller
#   timeout fires SIGKILL prematurely and noisy "StopSignal SIGTERM
#   failed to stop container" warnings flood the journal on every
#   restart cascade (template edit -> all instances bounce). 60s
#   covers the long-poll plus a little graceful-exit wiggle. Kept
#   strictly less than TimeoutStopSec above so podman wins the race
#   and systemd doesn't SIGKILL podman mid-cleanup.
# --pull=newer: do a manifest check against nexus on every container
#   start and pull only if the registry's image is newer than the one
#   in the local store. Lets a fresh runner-image.yml build propagate
#   to every instance on its next restart -- no ansible run needed.
#   Cheap when up-to-date (just a HEAD on the manifest, ~100ms);
#   first-after-bump start eats the layer download once, then the
#   layer cache satisfies subsequent starts until the next bump.
ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
  --cidfile=%t/%n.ctr-id \
  --cgroups=split \
  --sdnotify=healthy \
  --detach \
  --replace \
  --rm \
  --pull=newer \
  --name github_runner_%i \
  --log-driver journald \
  --log-opt tag="github_runner_%i" \
  --user 0:0 \
  --uidmap=0:0:65536 \
  --uidmap=+0:{{ github_runner_user.uid }}:1 \
  --gidmap=0:0:65536 \
  --gidmap=+0:{{ github_runner_user.group }}:1 \
  --volume %t/podman/podman.sock:/var/run/docker.sock \
  --volume /mnt/scratch/github_runner/workdir_%i:/work:rw \
  --env-file /etc/default/github_runner@%i \
  --secret github_runner_pat,type=mount,target=/run/secrets/github_pat \
  --health-cmd         "pgrep -f Runner.Listener >/dev/null" \
  --health-startup-cmd "pgrep -f Runner.Listener >/dev/null" \
  --health-startup-interval "1s" \
  --health-startup-timeout "5s" \
  --health-start-period "30s" \
  --stop-timeout 60 \
  nexus.lab.fahm.fr/homelab/runner:latest
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=default.target
{% if zfs_root %}
# Intent-documentation only (see [Unit] zfs comment above): the actual
# zfs_autosnapshot.target integration lives in the user@<uid>.service
# drop-in. Recorded here so the unit's autosnap relationship is
# legible at first read.
WantedBy=zfs_autosnapshot.target
{% endif %}
