# roles/github_runner/templates/github_runner.service.j2
# jinja2: lstrip_blocks: "True"
# github_runner@.service ###############################################
# Templated unit; one instance per entry in github_runner_instances on
# this host. Instance name = <repo>_<suffix> (e.g. homelab_vm1). Per-
# instance config (REPO_URL, RUNNER_NAME, RUNNER_LABELS) ships via
# /etc/default/github_runner@%i; first-start registration token ships
# via the github_runner_token_%i podman secret. State persists on the
# /mnt/services/github_runner/%i bind-mount so container restarts and
# image bumps don't re-register.
[Unit]
Description=github-actions runner (%i)
After=network-online.target
Wants=network-online.target
# actions/runner inside the container spawns sibling workflow containers
# via the bind-mounted rootless podman socket; wait for github_runner's
# user manager (which owns that socket) to be up before starting.
After=user@{{ github_runner_user.uid }}.service
Requires=user@{{ github_runner_user.uid }}.service
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
{% endif %}
{% if zfs_has_scratch_mount %}
After=zfs_mount_mnt_scratch.service
Requires=zfs_mount_mnt_scratch.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
RestartSec=5
Delegate=yes
# Cold-start has to fit: image pull (first time), container create, the
# entrypoint's rsync from /runner-base to /runner, and on first start
# also the actions/runner config.sh registration call. 300s covers the
# slowest of these (network-bound pull) with headroom.
TimeoutStartSec=300
TimeoutStopSec=60
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

User=github_runner
Group=github_runner

# Wait for the github_runner user manager to bring its rootless
# podman.socket up. Bounded 60 x 1s, well inside TimeoutStartSec. If the
# socket never appears, the journal records the explicit cause rather
# than a generic start-timeout.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -S /run/user/{{ github_runner_user.uid }}/podman/podman.sock ] && exit 0; sleep 1; done; echo "podman.sock not present after 60s" >&2; exit 1'

# --user / --uidmap / --gidmap: fake-root pattern (see CLAUDE.md
#   "Container User Namespacing"). In-container uid 0 maps to host
#   github_runner via the +0 override; everything else falls through
#   1:1 from the base 65536-range mapping. Writes to the bind-mounted
#   /runner state dir and the bound podman.sock land owned by host
#   github_runner. NOT `--userns keep-id`: that would map the *current
#   rootless user* to the same in-container uid, leaving in-container
#   root pointing at an unmapped subuid and failing every write to
#   /runner.
# --volume podman.sock: DooD -- actions/runner's docker-CLI shell-outs
#   (for workflows with `container:` blocks) hit the host's rootless
#   podman, NOT a daemon inside the runner container.
# --volume /runner: persisted runner state (binaries, .runner,
#   .credentials, _diag). Survives container restarts and image bumps.
# --volume /work: actions/runner's checkout root (_work tree). Lives on
#   /mnt/scratch (terabytes) so per-job checkouts don't pressure the
#   container's writable layer.
# --env-file: per-instance REPO_URL / RUNNER_NAME / RUNNER_LABELS,
#   written by the role to /etc/default/github_runner@%i.
# --secret RUNNER_TOKEN: short-lived registration token (1h TTL). Only
#   consumed on first start when /runner/.runner is absent; subsequent
#   starts ignore it and use the persisted .credentials.
# --health-cmd / startup-cmd: actions/runner exposes no HTTP probe, so
#   we go with CLAUDE.md "Container Healthchecks" option 3 (service-
#   native CLI): pgrep the Runner.Listener process. procps ships in
#   ubuntu:24.04's base so the runner image inherits it without an
#   explicit apt install. When the listener loses its GitHub
#   connection it tries to reconnect; if reconnection fails it exits
#   non-zero, the container exits, and systemd's Restart=always
#   respawns -- so the healthcheck is the "is it alive" signal
#   surfaced in `systemctl status`/downstream monitoring, while the
#   exit+respawn loop drives the actual recovery.
# --health-start-period 30s: first-start runs config.sh (network call
#   to GitHub for registration) before run.sh. 30s leaves headroom for
#   a slow API response without flapping the unit unhealthy on boot.
ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
  --cidfile=%t/%n.ctr-id \
  --cgroups=split \
  --sdnotify=conmon \
  --detach \
  --replace \
  --rm \
  --name github_runner_%i \
  --log-driver journald \
  --log-opt tag="github_runner_%i" \
  --user 0:0 \
  --uidmap=0:0:65536 \
  --uidmap=+0:{{ github_runner_user.uid }}:1 \
  --gidmap=0:0:65536 \
  --gidmap=+0:{{ github_runner_user.group }}:1 \
  --volume /run/user/{{ github_runner_user.uid }}/podman/podman.sock:/var/run/docker.sock \
  --volume /mnt/services/github_runner/%i:/runner:rw \
  --volume /mnt/scratch/github_runner/workdir_%i:/work:rw \
  --env-file /etc/default/github_runner@%i \
  --secret github_runner_token_%i,type=env,target=RUNNER_TOKEN \
  --health-cmd         "pgrep -f Runner.Listener >/dev/null" \
  --health-startup-cmd "pgrep -f Runner.Listener >/dev/null" \
  --health-startup-interval "1s" \
  --health-startup-timeout "5s" \
  --health-start-period "30s" \
  --stop-timeout 30 \
  nexus.lab.fahm.fr/homelab/runner:latest
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=default.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
