# roles/github_runner/templates/github_runner.service.j2
#jinja2: lstrip_blocks: True
# github_runner@.service ###############################################
# Templated system-scope unit running actions/runner directly on the
# host as the github_runner system user. No containerisation: each
# instance's run.sh executes in /opt/actions-runner/<inst>/ (a hard-
# link mirror of the canonical /opt/actions-runner/ binary tree plus
# the instance's own .runner / .credentials* config). Long-lived
# registration (config.sh runs once at apply time per instance via
# _register.yml); the unit just spawns run.sh and lets it long-poll
# GitHub.
#
# Previous shape (commit 82a1ac86 .. 286a5809) ran the runner inside a
# rootless-but-system-scope podman container with a JIT-config minted
# per container start, +0 uidmap demotion, DooD via bind-mounted
# rootless podman socket, sync_root hard-linking per-instance binary
# trees on each restart, etc. The whole edifice existed for security
# properties (per-job clean state, ephemeral credentials, namespace
# isolation) that matter at $BigCo and barely matter at single-
# operator homelab scale. Reverted to "actions/runner on the host"
# because the operational tax of the containerised shape (JIT mint
# script, per-start sync_root, image rebuild pipeline, uidmap+subgid
# acrobatics, /var/run/docker.sock symlink dance, last-write-wins
# races at /opt/actions-runner/) exceeded the threat-model benefit.
# A malicious workflow still gets capped at the github_runner uid,
# which has no privileged access to anything on the host.
[Unit]
Description=github-actions runner (%i)
After=network-online.target
Wants=network-online.target
# nexus is contacted by `container:` workflow blocks that pull from
# nexus.lab.fahm.fr/homelab/ci:latest. After= orders cold start so the
# first job's image pull doesn't race nexus bringup; Wants= (not
# Requires=) keeps a nexus stop from cascading into in-flight jobs.
After=nexus.service
Wants=nexus.service
# user@<uid>.service brings up github_runner's rootless podman.socket
# via linger. Workflow `container:` blocks shell out to docker (via
# DOCKER_HOST below) which is that rootless podman socket -- so the
# user manager has to be up before run.sh can accept jobs that need
# it.
After=user@{{ github_runner_user.uid }}.service
Requires=user@{{ github_runner_user.uid }}.service
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_scratch_mount %}
After=zfs_mount_mnt_scratch.service
Requires=zfs_mount_mnt_scratch.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Type=simple
User=github_runner
Group=github_runner
WorkingDirectory=/opt/actions-runner
# Wrapper script cd's to the per-instance dir and execs run.sh. The
# indirection through a stable /usr/local/bin/ path is what makes the
# unit pass systemd-analyze verify -- a direct
# ExecStart=/opt/actions-runner/%i/run.sh fails validation because
# verify substitutes "i" as the instance placeholder and the dir
# doesn't exist for that literal. The wrapper also handles
# WorkingDirectory (cd inside the script) so systemd doesn't try to
# chdir into a possibly-missing dir before forking ExecStart.
ExecStart=/usr/local/bin/github_runner_run %i
Restart=always
RestartSec=10
# actions/runner's main loop is a long-poll against GitHub's
# /runners/messages endpoint (~50s server hold). SIGTERM only takes
# effect once the in-flight long-poll returns, so 90s gives the listener
# room to acknowledge the signal before systemd escalates to SIGKILL.
TimeoutStopSec=90
KillSignal=SIGTERM
# Workflow `container:` blocks shell out to `docker` (actually
# podman-docker on the host) which connects to DOCKER_HOST. Point at
# github_runner's rootless podman socket so sibling job containers
# run under the same uid as the runner process. The /var/run/docker.sock
# symlink (see main.yml's tmpfiles + immediate-apply tasks) covers the
# rare workflow that hardcodes that path.
Environment=DOCKER_HOST=unix:///run/user/{{ github_runner_user.uid }}/podman/podman.sock
# actions/runner reads its work folder from .runner (set at config.sh
# time via --work). _register.yml points --work at
# /mnt/scratch/github_runner/workdir_<inst> so per-job checkouts +
# tool installs land on /mnt/scratch (terabytes) instead of /opt
# (a rpool ZFS dataset).
#
# Workflow-spawned `container:` blocks run under github_runner's
# rootless podman, so they inherit the same kvm-group access via the
# subid passthrough configured in subid_extra_subgid_groups (see
# host_vars/<host>.yml). Job-level options like
# `--device /dev/kvm --group-add 1` continue to work -- the "1" is
# the container-side gid that subgid maps to host kvm.

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
