# roles/healthchecks/templates/healthchecks.service.j2
# jinja2: lstrip_blocks: "True"
# healthchecks.service #######################################################################
[Unit]
Description=healthchecks
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=60
RestartSec=5
Type=notify
# Has to be `all`, not `main`. With --detach below podman exits as soon as
# the container is handed off to conmon, and conmon — not podman — is the
# process that watches the healthcheck and sends READY=1. conmon is never
# MainPID, so NotifyAccess=main causes the unit to hang in `activating`.
NotifyAccess=all
SyslogIdentifier=%N

# Both SUPERUSER_PASSWORD (used once at first-boot by linuxserver's
# s6 init) and SECRET_KEY (read every request by Django) flow in via
# --secret=type=env rather than type=mount/_FILE because
# linuxserver/healthchecks' s6-overlay init reads them straight from
# the env -- there's no _FILE-suffix handling in that wrapper. The
# env-var exposure (podman inspect, /proc/<pid>/environ) is
# 0400 root-only in this single-operator host, so the leak isn't
# reachable in the threat model. Reconsider if we move to an image
# whose wrapper does honor _FILE.
ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name healthchecks \
    --log-driver journald \
    --log-opt tag="healthchecks" \
    --security-opt=no-new-privileges \
    --uidmap=0:0:65536 \
    --uidmap=+0:{{ healthchecks_user.uid }}:1 \
    --gidmap=0:0:65536 \
    --gidmap=+0:{{ healthchecks_user.group }}:1 \
    --volume /mnt/services/healthchecks:/config \
    --volume /mnt/services/sqlite/healthchecks:/db \
    --volume /mnt/services/healthchecks/logo.png:/app/healthchecks/static-collected/img/logo.png:ro \
    --env PUID=0 \
    --env PGID=0 \
    --env DB_NAME=/db/hc.sqlite \
    --env RP_ID=healthchecks.{{ inventory_hostname }}.{{ domain }} \
    --env SITE_ROOT=https://healthchecks.{{ inventory_hostname }}.{{ domain }} \
    --env ALLOWED_HOSTS=healthchecks.{{ inventory_hostname }}.{{ domain }} \
    --env REGISTRATION_OPEN=False \
    --env SECURE_PROXY_SSL_HEADER="X-Forwarded-Proto,https" \
    --env CSRF_TRUSTED_ORIGINS=https://healthchecks.{{ inventory_hostname }}.{{ domain }} \
    --env TZ={{ system_timezone }} \
    --env SITE_NAME=Healthchecks \
    --env ADMINS={{ root_email }} \
    --env SUPERUSER_EMAIL={{ root_email }} \
    --secret=healthchecks_password,type=env,target=SUPERUSER_PASSWORD \
    --secret=healthchecks_secret_key,type=env,target=SECRET_KEY \
    --env DEBUG=False \
    --env DEFAULT_FROM_EMAIL="{{ healthchecks_smtp_from }}" \
    --env EMAIL_HOST="{{ smtp_host }}" \
    --env EMAIL_PORT="{{ smtp_port }}" \
    --env EMAIL_HOST_USER="{{ healthchecks_smtp_user}}" \
    --secret=healthchecks_smtp_password,type=mount,target=email_host_password \
    --env EMAIL_HOST_PASSWORD_FILE=/run/secrets/email_host_password \
    --env EMAIL_USE_TLS="{{ smtp_use_tls | bool }}" \
    --env EMAIL_USE_SSL="{{ smtp_use_ssl | bool }}" \
    --env SITE_LOGO_URL=/static/img/logo.png \
    --env UWSGI_DISABLE_LOGGING=1 \
    --env UWSGI_HOOK_MASTER_START="unix_signal:15 gracefully_kill_them_all" \
    --publish 127.0.0.1:{{ service_ports.healthchecks }}:8000/tcp \
    --health-cmd         '["curl","--head","--request","GET","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","--header","Host: healthchecks.{{ inventory_hostname }}.{{ domain }}","-o","/dev/null","http://127.0.0.1:8000/api/v3/status/"]' \
    --health-startup-cmd '["curl","--head","--request","GET","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","--header","Host: healthchecks.{{ inventory_hostname }}.{{ domain }}","-o","/dev/null","http://127.0.0.1:8000/api/v3/status/"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "10s" \
    --stop-timeout 30 \
    docker.io/linuxserver/healthchecks:4.2.20260428
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=default.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
