# roles/kuma/templates/kuma.service.j2
#jinja2: lstrip_blocks: True
# kuma.service #######################################################################
[Unit]
Description=kuma
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
# Start budget is 120s -- kuma's first-boot Database.connect() + Express
# wire-up has been observed in the 30-60s range; 120s leaves headroom.
# TimeoutStopSec sits above --stop-timeout 120 below, so systemd waits
# for podman's graceful-stop budget plus 30s slack before SIGKILLing.
# Without that slack, systemd's SIGKILL races podman's SIGKILL and can
# truncate the SQLite WAL flush mid-write.
TimeoutStartSec=120
TimeoutStopSec=150
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name kuma \
    --user {{ kuma_user.uid }}:{{ kuma_user.group }} \
    {# No --entrypoint override: the v2 release image runs
       `dumb-init -- node server/server.js` natively, with no chown-doing
       extra/entrypoint.sh wrapper to bypass. #}
    --log-driver journald \
    --log-opt tag="kuma" \
    --security-opt=no-new-privileges \
    {# v2 added a first-boot "choose your database" step: with neither a
       db-config.json nor an existing kuma.db present, the server blocks on a
       setup page and never mounts its real routes (api-router, /metrics) --
       yet extra/healthcheck still answers 200, so the unit goes `active` on a
       non-functional instance (silent on fresh hosts / disaster recovery).
       Pinning the type skips the setup page and writes db-config.json itself.
       (lab's existing kuma.db would auto-select sqlite, but new hosts wouldn't.) #}
    --env UPTIME_KUMA_DB_TYPE=sqlite \
    {# Pin the data dir: the kuma.db symlink and the /app/data volume below
       both assume kuma writes its DB at /app/data root. Declaring it makes
       that load-bearing assumption explicit rather than an image default. #}
    --env DATA_DIR=/app/data \
    --volume /mnt/services/kuma:/app/data \
    --volume /mnt/services/sqlite/kuma:/mnt/services/sqlite/kuma \
    --publish 127.0.0.1:{{ service_ports.kuma }}:3001/tcp \
    {# extra/healthcheck is kuma's own image-shipped liveness script (node),
       chosen over a hand-rolled tier-2 curl probe; UPTIME_KUMA_DB_TYPE above
       prevents the setup-page false-positive it would otherwise return. #}
    --health-cmd '["extra/healthcheck"]' \
    --health-startup-cmd '["extra/healthcheck"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 120 \
    {{ kuma_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
