# roles/gitea/templates/gitea.service.j2
#jinja2: lstrip_blocks: True
# gitea.service #######################################################################
[Unit]
Description=gitea
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Restart=always
# gitea cold start is <5s; the ExecStartPost curl probe below has a
# ~30s budget. 60s leaves headroom without doubling the unavailable
# window on a wedged binary.
TimeoutStartSec=60
TimeoutStopSec=30
RestartSec=5
SyslogIdentifier=%N

# gitea 1.20+ calls sd_notify after the HTTP listener is bound
# (graceful/manager_unix.go: notify(readyMsg)). NotifyAccess=all because
# gitea forks worker subprocesses; main alone is too narrow.
Type=notify
NotifyAccess=all
User={{ git_user.name }}
Group={{ git_user.name }}
WorkingDirectory=/mnt/services/gitea/
ExecStart=/usr/local/bin/gitea web --config /etc/gitea/app.ini
Environment=USER={{ git_user.name }} HOME={{ git_user.home }} GITEA_WORK_DIR=/mnt/services/gitea

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
