# roles/gitea/templates/gitea.service.j2
# jinja2: lstrip_blocks: "True"
# gitea.service #######################################################################
[Unit]
Description=gitea
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Restart=always
TimeoutStartSec=120
TimeoutStopSec=30
RestartSec=5
SyslogIdentifier=%N

Type=simple
User={{ git_user.name }}
Group={{ git_user.name }}
WorkingDirectory=/mnt/services/gitea/
ExecStart=/usr/local/bin/gitea web --config /etc/gitea/app.ini
# Block "active" until gitea is actually answering HTTP. --retry-connrefused
# bounds at ~30s; ExecStartPost is the safe alternative to Type=notify (gitea
# may or may not call sd_notify depending on build).
ExecStartPost=/usr/bin/curl --fail --silent --max-time 30 --retry 30 --retry-delay 1 --retry-connrefused http://127.0.0.1:{{ service_ports.gitea }}/api/healthz
Environment=USER={{ git_user.name }} HOME={{ git_user.home }} GITEA_WORK_DIR=/mnt/services/gitea

[Install]
WantedBy=default.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
