# roles/sort_ini/templates/sort_ini@.path.j2
[Unit]
Description="Monitor {{ path }} file for changes"

[Path]
PathChanged={{ path }}

[Install]
WantedBy=multi-user.target
