# roles/sort_ini/templates/sort_ini.path.j2
[Unit]
Description="Monitor {{ sort_ini_path }} file for changes"

[Path]
PathModified={{ sort_ini_path }}

[Install]
WantedBy=multi-user.target
