# roles/compta/templates/compta.service.j2
#jinja2: lstrip_blocks: True
# compta.service #######################################################################
[Unit]
Description=compta
After=network-online.target
Wants=network-online.target
# Image lives at nexus.lab.fahm.fr/compta/compta:latest; nexus is co-located
# on this host (site.yml: same `hosts: box,lab` block). Wants= (not Requires=)
# so a nexus stop doesn't cascade-stop compta — the cached image still runs.
After=nexus.service
Wants=nexus.service
StartLimitInterval=600
StartLimitBurst=3

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=30
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --init \
    --name compta \
    --log-driver journald \
    --log-opt tag="compta" \
    --security-opt=no-new-privileges \
    --user {{ compta_user.uid }}:{{ compta_user.group }}  \
    --publish 127.0.0.1:{{ service_ports.compta }}:4567/tcp \
    --health-cmd         '["curl","--head","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:4567/"]' \
    --health-startup-cmd '["curl","--head","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:4567/"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "15s" \
    --stop-timeout 10 \
    nexus.lab.fahm.fr/compta/compta:latest
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
