# roles/gitlab_runner/templates/gitlab_runner.service.j2
#jinja2: lstrip_blocks: True
# gitlab_runner.service ###############################################
# GitLab Runner, host-native: a single static Go binary at
# /usr/local/bin/gitlab-runner running as the gitlab_runner system user.
# Docker executor against the user's rootless podman socket -- this
# process runs job payloads in containers it starts over that socket,
# never directly as itself. The runner process itself is not containerised.
[Unit]
Description=GitLab Runner (Docker executor -> rootless podman)
After=network-online.target
Wants=network-online.target
# Restart=always + RestartSec=10 can trip systemd's default 5-in-10s burst
# during a slow network warmup and land the unit in `failed`; widen the window.
StartLimitIntervalSec=600
StartLimitBurst=5
{% if zfs_root %}
{% if zfs_has_services_mount %}
# config.toml + the runner's working dir live under /mnt/services; order after
# the dataset mount so they exist before start.
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
{% endif %}
{% endif %}
# Deliberately NOT PartOf=zfs_autosnapshot.target: the runner has no persistent
# consistency need (config.toml is ansible-rendered; job state lives in
# throwaway containers), and joining the target would let the nightly 02:00
# snapshot stop the runner mid-CI.

[Service]
Type=simple
User=gitlab_runner
Group=gitlab_runner
WorkingDirectory=/mnt/services/gitlab_runner
# nologin service user has no real home; point HOME at the working dir so any
# runner cache lands on a writable, owned path rather than failing against an
# unset/unwritable HOME.
Environment=HOME=/mnt/services/gitlab_runner
ExecStart=/usr/local/bin/gitlab-runner run \
    --config /mnt/services/gitlab_runner/config.toml \
    --working-directory /mnt/services/gitlab_runner
Restart=always
RestartSec=10
# gitlab-runner drains gracefully on SIGQUIT: it stops claiming new jobs and
# waits for in-flight jobs to finish, bounded by config.toml's
# shutdown_timeout (1800s), then exits. TimeoutStopSec sits just past that
# bound as the SIGKILL backstop: long enough that systemd never truncates a
# healthy drain, short enough that a runner too wedged to honor SIGQUIT
# cannot block stop/restart -- and with them converge and host reboot --
# forever.
KillSignal=SIGQUIT
TimeoutStopSec=1860

[Install]
WantedBy=multi-user.target
