# roles/github_runner/files/buildkit.service
[Unit]
Description=Rootless BuildKit for github_runner CI image builds
After=podman.socket
# Wants (not BindsTo): buildkitd only uses the podman socket when `docker
# buildx --load` fires -- not continuously. BindsTo would co-stop buildkitd
# whenever podman.socket restarts (e.g. during a podman role converge), killing
# any in-progress zbm-build CI job.
Wants=podman.socket

[Service]
Environment=PATH=/usr/local/bin:/usr/bin:/bin
RuntimeDirectory=buildkit
# --net=host: rootlesskit defaults to slirp4netns NAT for outbound traffic. buildkitd
# needs no network isolation (it is a build executor, not an untrusted container),
# so host networking avoids the NAT overhead on every xbps package download.
# --oci-worker-gc-keepstorage: minimum retained storage in bytes. The ZBM
# builder image alone is ~1.3 GB; 20 GiB retains it plus the xbps download cache across
# version bumps (the runner's /mnt/scratch has ~1 TB). Without this, buildkitd's 2 GB
# default evicts the image between runs and forces full cold rebuilds.
# TimeoutStopSec: buildkitd flushes its cache to disk on shutdown; the systemd default
# (90 s) may not be enough for a 20 GB cache, and a SIGKILL-interrupted flush can corrupt
# the xbps layer cache that the mount-type-cache steps populate.
ExecStart=/usr/local/bin/rootlesskit --net=host /usr/local/bin/buildkitd --addr unix://%t/buildkit/buildkitd.sock --root %h/.local/share/buildkit --oci-worker-gc-keepstorage 21474836480
Restart=on-failure
RestartSec=5s
TimeoutStopSec=300

[Install]
WantedBy=default.target
