# roles/mosquitto/templates/mosquitto.service.j2
#jinja2: lstrip_blocks: True
# mosquitto.service #######################################################################
[Unit]
Description=mosquitto
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=120
RestartSec=5
Type=notify
# NotifyAccess=all (not main): with --sdnotify=healthy the readiness ping comes
# from podman/conmon relaying the container healthcheck, not the unit's main
# PID -- tightening to `main` would stop the unit ever reaching active.
NotifyAccess=all
SyslogIdentifier=%N

# The healthcheck probes the loopback 1880 listener, which mosquitto.conf
# keeps anonymous so mosquitto_sub needs no credentials; -E exits on the
# broker's subscribe-ack, which is the liveness signal. Ad-hoc debugging is
# via `podman exec mosquitto mosquitto_sub ...`, not a published host port.
ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name mosquitto \
    --init \
    --log-driver journald \
    --log-opt tag="mosquitto" \
    --volume /mnt/services/mosquitto/config:/mosquitto/config \
    --volume /mnt/services/sqlite/mosquitto:/mosquitto/data \
    --user {{ mosquitto_user.uid }}:{{ mosquitto_user.group }} \
    --network homeassistant \
    --health-cmd {{ ["mosquitto_sub", "-p", "1880", "-t", "healthcheck", "-E", "-i", "probe", "-W", "3"] | json_argv }} \
    --health-startup-cmd {{ ["mosquitto_sub", "-p", "1880", "-t", "healthcheck", "-E", "-i", "probe", "-W", "3"] | json_argv }} \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 120 \
    {{ mosquitto_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
