# .old/ansible/roles/services/files/systemd/docker-hub-schedule.timer
[Timer]
OnBootSec=15min
OnUnitActiveSec=1w

[Install]
WantedBy=compose.service
