# .old/ansible/roles/services/files/systemd/pi-backup.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/dms 98d0602056 /usr/local/bin/pi-backup
