# .old/ansible/roles/services/files/systemd/usenet-maintenance.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/usenet-maintenance
