# roles/jellyfin/templates/jellyfin.service.j2
#jinja2: lstrip_blocks: True, trim_blocks: True
# jellyfin.service #######################################################################
[Unit]
Description=jellyfin
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
{% if zfs_has_scratch_mount %}
After=zfs_mount_mnt_scratch.service
Requires=zfs_mount_mnt_scratch.service
{% endif %}
{% if zfs_has_media_mount %}
After=zfs_mount_mnt_media.service
Requires=zfs_mount_mnt_media.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
# First-boot EF Core migrations on the SQLite library DB take >100s
# on emulated aarch64 (the default 100s health-start-period below
# would otherwise time the unit out before sdnotify=healthy fires).
# TimeoutSec must outlast --health-start-period for the wait to win.
TimeoutSec=360
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name jellyfin \
    --init \
    --user 0:0 \
    --uidmap=0:0:65536 \
    --uidmap=+0:{{ jellyfin_user.uid }}:1 \
    --gidmap=0:0:65536 \
    --gidmap=+0:{{ jellyfin_user.group }}:1 \
    --group-add={{ media_group.gid }} \
    --log-driver journald \
    --log-opt tag="jellyfin" \
    --volume /mnt/services/jellyfin:/config \
    --volume /mnt/services/sqlite/jellyfin:/mnt/services/sqlite/jellyfin \
    --mount type=tmpfs,destination=/config/transcodes,tmpfs-size=4G,tmpfs-mode=0755 \
    --volume /var/log/jellyfin:/config/log \
    --volume /mnt/scratch/jellyfin/cache:/cache \
    --volume /mnt/media:/media \
{% if has_transcoding_acceleration %}
    --device /dev/dri/:/dev/dri/ \
    --group-add={{ render_group.gid }} \
{% endif %}
    --env JELLYFIN_CACHE_DIR=/cache \
    --env JELLYFIN_PublishedServerUrl=https://jellyfin.{{ inventory_hostname }}.{{ domain }} \
    --network media \
    --publish 127.0.0.1:{{ service_ports.jellyfin }}:8096/tcp \
    --health-cmd         '["curl","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:8096/health"]' \
    --health-startup-cmd '["curl","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:8096/health"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "300s" \
    --stop-timeout 120 \
    docker.io/jellyfin/jellyfin:10.10.6
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
