# roles/webmail/templates/roundcube.service.j2
#jinja2: lstrip_blocks: True
# roundcube.service #######################################################################
# Webmail UI for the getmail archive. Talks IMAP to the dovecot container over
# the internal webmail network; nginx fronts it with Authelia auth.
[Unit]
Description=roundcube
After=network-online.target
Wants=network-online.target
# dovecot.service is Type=notify --sdnotify=healthy, so ordering after it means
# Roundcube boots against a ready IMAP backend (and guarantees the shared
# webmail network exists). Wants, not Requires: Roundcube tolerates
# dovecot bouncing — the user just re-logs in.
After=dovecot.service
Wants=dovecot.service
StartLimitIntervalSec=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutStartSec=120
TimeoutStopSec=40
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name roundcube \
    --log-driver journald \
    --log-opt tag="roundcube" \
    --security-opt=no-new-privileges \
    --network webmail \
    {# The webroot is ephemeral: the image re-extracts its whole tree into
       /var/www/html on every boot (--rm wipes the prior copy), so des_key and
       the env-rendered config.inc.php regenerate each start and a restart
       forces a re-login — fine for an occasional archive viewer. The extract
       chmods every dir via fchmodat2; podman's default seccomp profile blocks
       that until the podman role's fchmodat2 allow-list lands (see
       notes/podman_seccomp_fchmodat2.md), without which apache 403s the lot. #}
    {# SQLite is the real state (prefs/contacts/identities) and persists on the
       dedicated 4K-recordsize sqlite dataset. ROUNDCUBEMAIL_DB_DIR is a
       dedicated mount point, so we point it straight at the dataset — no
       symlink dance (unlike kuma/paperless, whose DB sits inside a larger data
       dir). The image's entrypoint creates sqlite.db here on first boot. #}
    --volume /mnt/services/sqlite/roundcube:/var/roundcube/db \
    {# Disable compose/reply/forward UI and SMTP — read-only archive viewer. #}
    --volume /mnt/services/webmail/roundcube_disable_send.php:/var/roundcube/config/disable-send.php:ro \
    --publish 127.0.0.1:{{ service_ports.roundcube }}:80/tcp \
    {# Plain IMAP on 143 to the dovecot container; the hop stays on the internal
       podman bridge so TLS would be ceremony. #}
    --env ROUNDCUBEMAIL_DEFAULT_HOST=dovecot.dns.podman \
    --env ROUNDCUBEMAIL_DEFAULT_PORT=143 \
    --env ROUNDCUBEMAIL_DB_DIR=/var/roundcube/db \
    --env ROUNDCUBEMAIL_DB_TYPE=sqlite \
    --env ROUNDCUBEMAIL_SKIN=elastic \
    {# Skip Roundcube's own login form: the http_authentication plugin logs in
       from the HTTP Basic credentials nginx injects (archive:archive — the same
       non-secret IMAP cred), so the only gate the user sees is Authelia. mod_php
       populates PHP_AUTH_USER/PW from the Authorization header natively. See the
       nginx location_conf in tasks/main.yml and notes/getmail_webmail_design.md. #}
    --env ROUNDCUBEMAIL_PLUGINS=http_authentication \
    --health-cmd         '["curl","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:80/"]' \
    --health-startup-cmd '["curl","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:80/"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "60s" \
    --stop-timeout 30 \
    docker.io/roundcube/roundcubemail:{{ webmail_roundcube_version }}-apache
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
