# roles/bazarr/templates/bazarr.service.j2
#jinja2: lstrip_blocks: True
# bazarr.service #######################################################################
[Unit]
Description=bazarr
After=network-online.target
Wants=network-online.target
# Order-only (no Requires=): on hosts where sonarr/radarr aren't deployed
# (test variants, future hostsplits) bazarr still starts; on pug it waits
# so the analytics tracker's import-time HTTPS calls land on a running
# sonarr/radarr instead of timing out 60s/call. See
# notes/bazarr_test_blocking_imports.md.
After=sonarr.service radarr.service
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
{% if zfs_has_media_mount %}
After=zfs_mount_mnt_media.service
Requires=zfs_mount_mnt_media.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=120
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name bazarr \
    --security-opt=no-new-privileges \
    --log-driver journald \
    --log-opt tag="bazarr" \
    {# PUID/PGID (CLAUDE.md "User namespacing" tier 2): the linuxserver s6
       init usermods its baked-in abc user to these ids and drops to it.
       Tier 2 is the right fit for a droppable image -- the fake-root
       --uidmap pattern is reserved for images that insist on id 0. #}
    --env PUID={{ bazarr_user.uid }} \
    --env PGID={{ media_group.gid }} \
    --env UMASK=002 \
    --env TZ={{ system_timezone }} \
    --volume /mnt/services/bazarr:/config \
    --volume /mnt/services/sqlite/bazarr:/config/db \
    --volume /mnt/media:/media \
    --network media \
    --publish 127.0.0.1:{{ service_ports.bazarr }}:{{ service_ports.bazarr }}/tcp \
    --health-cmd         '["curl","--head","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:{{ service_ports.bazarr }}/api/"]' \
    --health-startup-cmd '["curl","--head","--location","--fail","--silent","--show-error","--connect-timeout","1","--max-time","5","-o","/dev/null","http://localhost:{{ service_ports.bazarr }}/api/"]' \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 60 \
    {{ bazarr_image }}
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
