# roles/bazarr/templates/bazarr.service.j2
# jinja2: lstrip_blocks: "True"
# bazarr.service #######################################################################
[Unit]
Description=bazarr
After=network-online.target
Wants=network-online.target
# Order-only (no Requires=): on hosts where sonarr/radarr aren't deployed
# (test variants, future hostsplits) bazarr still starts; on pug it waits
# so the analytics tracker's import-time HTTPS calls land on a running
# sonarr/radarr instead of timing out 60s/call. See
# notes/bazarr_test_blocking_imports.md.
After=sonarr.service radarr.service
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
{% if zfs_has_media_mount %}
After=zfs_mount_mnt_media.service
Requires=zfs_mount_mnt_media.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=120
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name bazarr \
    --security-opt=no-new-privileges \
    --log-driver journald \
    --log-opt tag="bazarr" \
    {# Still on PUID/PGID rather than the canonical --user 0:0 + --uidmap=+0:N:1
       pattern (CLAUDE.md "Container User Namespacing"). The linuxserver/bazarr
       image ships /dev/console as a char device in its layer; podman's
       chown-by-maps pass returns EIO trying to lchown it on the overlay
       driver, so the container won't start with uidmaps. Healthchecks's image
       doesn't have that device node so its uidmap migration works. The real
       fix is fleet-wide ignore_chown_errors=true in podman storage.conf;
       defer until that lands. #}
    --env PUID={{ bazarr_user.uid }} \
    --env PGID={{ media_group.gid }} \
    --env UMASK=002 \
    --env TZ={{ system_timezone }} \
    --volume /mnt/services/bazarr:/config \
    --volume /mnt/services/sqlite/bazarr:/config/db \
    --volume /mnt/media:/media \
    --publish 127.0.0.1:{{ service_ports.bazarr }}:{{ service_ports.bazarr }}/tcp \
    --health-cmd         "curl --head --location --fail --silent --show-error --connect-timeout 1 --max-time 5 -o /dev/null http://localhost:{{ service_ports.bazarr }}/api/" \
    --health-startup-cmd "curl --head --location --fail --silent --show-error --connect-timeout 1 --max-time 5 -o /dev/null http://localhost:{{ service_ports.bazarr }}/api/" \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 60 \
    docker.io/linuxserver/bazarr:1.5.2
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=default.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
