# roles/headscale/templates/headscale.service.j2
#jinja2: lstrip_blocks: True
# headscale.service ##################################################################
[Unit]
Description=headscale
After=network-online.target
Wants=network-online.target
# Baseline start-limit, present regardless of OIDC: Restart=always + RestartSec=5
# can trip systemd's default 5-starts-in-10s burst during a slow dependency
# warmup and land the unit in `failed` rather than retrying. Matches the sibling
# podman units (e.g. roles/kuma). OIDC widens the burst to 15 -- fox is a
# wg *listener*, so the tunnel to lab only carries traffic once lab/pug send their
# next keepalive (~25s after a fox reboot), and headscale must keep retrying
# boot-time OIDC discovery through that wg0 warmup rather than exhausting the
# burst. A genuinely-broken server still surfaces (Restart=always loops, and
# fox's netdata headscale collector goes missing -- see host_vars/fox.yml).
StartLimitIntervalSec=600
StartLimitBurst={{ 15 if _headscale_oidc_configured else 3 }}
{% if _headscale_oidc_configured %}
{% if not qemu_test %}
# OIDC discovery hits the Authelia issuer at startup, and on fox that issuer
# (auth.fahm.fr) resolves to lab's wg /32 -- reachable only over the wg0
# backbone. Order after wg-quick so the interface + route exist before
# headscale's first start attempt. fox always runs the wireguard role, and this
# unit only ever runs on fox, so the dep always holds. Skipped under qemu_test:
# the OIDC boot test points at a local mockoidc issuer (no wg0 in the harness),
# and an After= on a unit that doesn't exist would warn under systemd-analyze.
After=wg-quick@wg0.service
{% endif %}
{% endif %}
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutStartSec=120
TimeoutStopSec=70
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name headscale \
    --network headscale \
    --user {{ headscale_user.uid }}:{{ headscale_user.group }} \
    --log-driver journald \
    --log-opt tag="headscale" \
    --security-opt=no-new-privileges \
    {# tini PID 1: reaps zombies and forwards signals -- notably the SIGHUP the
       unit's ExecReload sends for an ACL policy reload. Matches sibling units. #}
    --init \
    --volume /mnt/services/headscale/config.yaml:/etc/headscale/config.yaml:ro \
    --volume /mnt/services/headscale/acl.hujson:/etc/headscale/acl.hujson:ro \
    --volume /mnt/services/headscale/data:/var/lib/headscale \
    {# db.sqlite under data/ is a symlink onto the dedicated 4K-recordsize
       sqlite dataset; mount that dataset at the same path inside the
       container so the absolute symlink target resolves (WAL/SHM siblings
       then also land on the dataset). See roles/sqlite_dataset. #}
    --volume /mnt/services/sqlite/headscale:/mnt/services/sqlite/headscale \
    --publish 127.0.0.1:{{ service_ports.headscale }}:8080/tcp \
    {# Prometheus metrics (config metrics_listen_addr 0.0.0.0:9090) on host
       loopback for fox's netdata go.d prometheus collector. Unauthenticated,
       so loopback-only -- never the public 3478/443 surface. #}
    --publish 127.0.0.1:{{ service_ports.headscale_metrics }}:9090/tcp \
    {# DERP STUN published on both families -- fox is a dual-stack relay. The
       container's STUN binds [::]:3478 (dual-stack) on the headscale podman
       network's ULA; netavark DNATs host v4 and v6 to it, source preserved. #}
    --publish 0.0.0.0:{{ service_ports.headscale_stun }}:3478/udp \
    --publish '[::]:{{ service_ports.headscale_stun }}:3478/udp' \
{% if _headscale_oidc_configured %}
    --secret=headscale_oidc_secret,type=mount,target=oidc_client_secret \
{% endif %}
    {# Liveness via the service-native CLI: `headscale users list` round-trips
       through the gRPC API over the unix socket to the running serve process,
       so it proves the server is actually answering (not just that the
       process exists). The headscale image is a Go-binary base with no curl,
       so the HTTP /health probe isn't available in-container. #}
    --health-cmd         {{ ["headscale", "users", "list"] | json_argv }} \
    --health-startup-cmd {{ ["headscale", "users", "list"] | json_argv }} \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "60s" \
    --stop-timeout 30 \
    {{ headscale_image }} \
    serve
# ACL-only changes reload via SIGHUP, not a restart: headscale reloads the
# HuJSON policy on HUP (same as `systemctl reload`), so a policy tweak needn't
# bounce the control plane + embedded DERP + every client long-poll. Fails soft
# -- a malformed policy leaves the running server on its previous ACL. Driven by
# the systemd_unit helper's reload: arg (roles/headscale/tasks/main.yml), which
# only reloads when the unit was already active. --init forwards the HUP to the
# headscale process.
ExecReload=/usr/bin/podman kill --signal HUP headscale
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
