# roles/firewall/files/wg0_flowtable.service
# Dynamically attaches wg0 to the inet homelab `bulk` flowtable. Kept
# out of /etc/nftables.conf's static `devices = { ... }` clause so
# nftables.service can load at boot without racing wg-quick@wg0 — nft
# errors out with ENOENT when a flowtable device doesn't exist yet,
# which would leave the host without a firewall until manual restart.
#
# PartOf=wg-quick@wg0.service nftables.service: a restart of either
# parent propagates a restart to us, so the dynamic attachment is
# re-applied after wg-quick comes back up or after a full
# nftables.service restart. PartOf doesn't propagate reload — the
# reload case is handled by nftables.service's ExecReload drop-in
# which calls `systemctl try-restart wg0_flowtable.service`.
#
# WantedBy=wg-quick@wg0.service: enabling installs the symlink under
# wg-quick@wg0.service.wants/, so starting wg-quick at boot also
# starts this unit. After= orders us after wg-quick has finished
# bringing wg0 up. Deliberately NOT BindsTo=wg-quick@wg0.service: a
# Requires-style dep would try to auto-start wg-quick whenever we're
# started, which under the test harness (where the wireguard role
# wasn't applied, so /etc/wireguard/wg0.conf doesn't exist) would
# cascade-fail the converge. PartOf + WantedBy gives the same
# lifecycle binding for prod (where wg-quick starts first) without
# the Requires-style coupling.
#
# ConditionPathExists silently skips the unit (no failed state, no
# emitted error in the journal) when wg0 doesn't exist. Covers the
# test-harness case where the wireguard role hasn't been applied and
# /sys/class/net/wg0 isn't present, plus a degraded-wg-quick boot
# where the iface never came up. PartOf will re-fire ExecStart
# whenever wg-quick is (re)started and brings wg0 up; ConditionPath
# is re-evaluated on each start.
#
# Not `ExecCondition=/bin/sh -c '... || exit 255'`: systemd docs say
# exit-255 *should* skip the unit silently, but in practice on jammy's
# systemd 249 the unit ends up in 'failed' state and `systemctl start
# --wait` returns non-zero. The static Condition* directives are the
# documented "silent skip" mechanism and behave consistently across
# systemd versions.
[Unit]
Description=Attach wg0 to nftables homelab flowtable
After=wg-quick@wg0.service nftables.service
PartOf=wg-quick@wg0.service nftables.service
ConditionPathExists=/sys/class/net/wg0

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/nft add flowtable inet homelab bulk { devices = { wg0 }; }

[Install]
WantedBy=wg-quick@wg0.service
