# roles/homeassistant/templates/homeassistant.service.j2
#jinja2: lstrip_blocks: True, trim_blocks: True
# homeassistant.service #######################################################################
[Unit]
Description=homeassistant
After=network-online.target
Wants=network-online.target
StartLimitInterval=600
StartLimitBurst=3
{% if zfs_root %}
{% if zfs_has_services_mount %}
After=zfs_mount_mnt_services.service
Requires=zfs_mount_mnt_services.service
After=zfs_mount_mnt_services_sqlite.service
Requires=zfs_mount_mnt_services_sqlite.service
{% endif %}
After=zfs_autosnapshot.target
PartOf=zfs_autosnapshot.target
{% endif %}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
Delegate=yes
TimeoutSec=120
RestartSec=5
Type=notify
NotifyAccess=all
SyslogIdentifier=%N

ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=split \
    --sdnotify=healthy \
    --detach \
    --replace \
    --rm \
    --name homeassistant \
    --log-driver journald \
    --log-opt tag="homeassistant" \
    --user root \
    {{ homeassistant_user | podman_idmap_args | join(" \\\n    ") }} \
{% for net in homeassistant_networks %}
{%   set _opts = [] %}
{%   if net.ip is defined %}{% set _ = _opts.append('ip=' ~ net.ip) %}{% endif %}
{%   if net.mac is defined %}{% set _ = _opts.append('mac=' ~ net.mac) %}{% endif %}
    --network={{ net.name }}{{ (':' ~ _opts | join(',')) if _opts else '' }} \
{% endfor %}
    --env TZ={{ system_timezone }} \
    --tz={{ system_timezone }} \
    --volume /mnt/services/homeassistant:/config \
    --volume /mnt/services/sqlite/homeassistant:/config/db \
    --publish 127.0.0.1:{{ service_ports.homeassistant }}:8123/tcp \
    {# Sonos S2 event-subscription callback; HomeKit pairing. Bound on the #}
    {# host's home-LAN address (not 0.0.0.0) so the iot-VLAN sub-iface #}
    {# doesn't host these sockets too -- iOS / Sonos controllers live on #}
    {# the home LAN. HA's macvlan presence on the IoT VLAN (10.123.4.0/24) #}
    {# remains reachable for mDNS / discovery; that's a separate plane. #}
    --publish {{ ansible_default_ipv4.address | mandatory }}:{{ service_ports.homeassistant_sonos }}:1400/tcp \
    --publish {{ ansible_default_ipv4.address | mandatory }}:{{ service_ports.homeassistant_homekit }}:3722/udp \
    --health-cmd         {{ "http://localhost:8123/manifest.json" | podman_health_curl }} \
    --health-startup-cmd {{ "http://localhost:8123/manifest.json" | podman_health_curl }} \
    --health-startup-interval "1s" \
    --health-startup-timeout "5s" \
    --health-start-period "100s" \
    --stop-timeout 120 \
    {{ homeassistant_image }}
{# Pin HA's default route once the netns is up. Where netavark installs a #}
{# default via podman bridge networks and the mac_iot macvlan at equal      #}
{# metric (an unpinned tie -- the shape lab is in), down-rank the macvlan  #}
{# one so off-IoT egress (home-VLAN media like the Marantz, WAN) takes a   #}
{# bridge NAT path -- sourced as the host's mgmt IP, which sidesteps the   #}
{# iot->* firewall rules. Where the macvlan is the *only* default (single- #}
{# default hosts, e.g. the box test fixture) the script leaves it as-is so #}
{# HA is never stranded; the mac_iot connected /24 route is untouched      #}
{# either way (in-VLAN IoT + mDNS preserved). We keep the network gateway  #}
{# (other containers may use mac_iot as sole egress), so this is scoped to #}
{# HA's netns. HA has no CAP_NET_ADMIN, so the script edits from the host  #}
{# via nsenter. Targeting is by gateway IP, not ethN -- podman's interface #}
{# ordering is not stable across hosts. $1 = mac_iot's gateway (its parent #}
{# subnet's .1, same derivation roles/macvlan uses). See files/pin_default_route.sh. #}
{% set _mac_iot_parent = (macvlan_blocks | one_by_attr('name', 'mac_iot')).parent %}
ExecStartPost=/opt/homeassistant/pin_default_route.sh {{ ansible_facts[_mac_iot_parent].ipv4.network | ansible.utils.ipmath(1) }}
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm --force --ignore --cidfile=%t/%n.ctr-id

[Install]
WantedBy=multi-user.target
{% if zfs_root %}
WantedBy=zfs_autosnapshot.target
{% endif %}
